# systemd/units/hapax-v4l2-bridge.service.d/darkplaces-runtime-gate.conf
[Unit]
# The bridge is owned by studio-compositor. When the DarkPlaces runtime marker
# is active, studio-compositor is intentionally condition-skipped; starting the
# bridge in that mode turns a valid runtime selection into repeated OnFailure
# P0 intake.
ConditionPathExists=!%h/.config/hapax/enable-darkplaces-runtime
