# systemd/units/codex-claim-audit.service.d/99-findings-are-warning.conf
[Service]
# codex-claim-audit exits 1 when unresolved claim/lane findings remain.
# The script already emits ntfy warnings; keep user-unit failed state for
# execution failures, not for an audit result with findings.
SuccessExitStatus=1
