# systemd/units/codex-claim-audit.service.d/99-findings-are-warning.conf
[Service]
# codex-claim-audit exits 10 when unresolved claim/lane findings remain — a
# warning, not an execution failure (the script already emits ntfy warnings).
# Treat ONLY that distinct code as success. A set -euo pipefail crash exits with
# the failing command's status (1/127/...) and release/reconcile failures exit 2;
# both must still leave the unit failed for this t1_critical claim-plane tool.
SuccessExitStatus=10
