# systemd/units/mediamtx.service.d/20-source-activation-worktree.conf
[Service]
WorkingDirectory=%h/.cache/hapax/source-activation/worktree
ExecStart=
ExecStart=%h/.cache/hapax/source-activation/worktree/scripts/mediamtx-start.sh
