# systemd/units/hapax-v4l2-bridge-watchdog.service.d/darkplaces-runtime-gate.conf
[Unit]
# Do not watchdog the studio-compositor V4L2 sidecar while the DarkPlaces
# runtime owns the visual path. The normal compositor service is condition-
# skipped in that mode, so watchdog restarts create synthetic bridge failures.
ConditionPathExists=!%h/.config/hapax/enable-darkplaces-runtime
