# systemd/units/hapax-broadcast-audio-health.service.d/90-source-activation-context.conf
[Service]
WorkingDirectory=%h/.cache/hapax/source-activation/worktree
ExecStart=
ExecStart=%h/.cache/hapax/source-activation/worktree/.venv/bin/python -m agents.broadcast_audio_health --loudness-duration 10
