# systemd/units/studio-compositor.service.d/cpu-affinity.conf
# Exclude studio-compositor from audio cores {6, 7, 14, 15}; allow {0-5, 8-13}.
# Pairs with the pipewire/wireplumber/pipewire-pulse pins to ensure compositor
# CPU spikes cannot preempt audio data-loops.
# See systemd/overrides/audio-stability/README.md for rationale.
[Service]
CPUAffinity=0 1 2 3 4 5 8 9 10 11 12 13
