# systemd/units/studio-compositor.service.d/layout-mode-persist.conf
[Service]
ExecStartPost=%h/.cache/hapax/source-activation/worktree/scripts/studio-compositor-post-start.sh
ExecStop=%h/.cache/hapax/source-activation/worktree/scripts/studio-compositor-persist-mode.sh
