# systemd/units/studio-compositor.service.d/malloc-arena.conf
[Service]
# Critical livestream incident 2026-05-13:
# studio-compositor uses many Cairo/GStreamer/Python worker threads. With the
# glibc default arena policy, live canary sampling showed ~10 MB/s RSS and
# Private_Dirty growth as repeated frame/surface allocations churned through
# per-thread arenas. Constrain arenas so allocator fragmentation is bounded
# while preserving every compositor visual capability.
Environment=MALLOC_ARENA_MAX=2
