# systemd/units/studio-compositor.service.d/oom-protect.conf
[Service]
# Start at the valid user-manager score, then let the root helper apply -800 to
# the authoritative service cgroup without protecting unrelated user children.
OOMScoreAdjust=100
ExecStartPost=-/usr/bin/sudo -n /usr/local/sbin/hapax-oom-score-enforce --apply-unit %n
MemoryLow=6G
MemoryMin=3G
