# systemd/units/studio-compositor.service.d/oom-protect.conf
[Service]
# Live broadcast compositor must not rely on stale inherited -900 scores.
OOMScoreAdjust=-800
MemoryLow=6G
MemoryMin=3G
