# systemd/units/studio-compositor.service.d/v4l2-bridge.conf
[Service]
# shmsink sidecar is the canonical v4l2 output path. The compositor writes
# frames to /dev/shm/hapax-compositor/v4l2-bridge.sock via shmsink; the
# hapax-v4l2-bridge sidecar reads from shmsrc and pushes to /dev/video42.
# This isolates the compositor process from v4l2loopback kernel state.
# Rollback: set HAPAX_V4L2_BRIDGE_ENABLED=0 to revert to direct os.write.
Environment=HAPAX_V4L2_BRIDGE_ENABLED=1
