# systemd/units/studio-compositor.service.d/v4l2-bridge.conf
[Service]
# Direct 3D V4L2 is the canonical live egress path while
# HAPAX_3D_COMPOSITOR=1 and hapax-imagination owns /dev/video42. Keep the
# legacy shmsink bridge sidecar opt-in; set this to 1 in a later runtime
# drop-in to prove bridge fallback mode without changing the bridge code.
Environment=HAPAX_V4L2_BRIDGE_ENABLED=0
