# systemd/units/studio-compositor.service.d/zzzz-screwm-quake-primary.conf
[Unit]
Wants=hapax-v4l2-bridge.service hapax-obs-video50-yuyv-compat-bridge.service
After=hapax-darkplaces-v4l2.service

[Service]
# Screwm Quake migration runtime profile, validated on 2026-05-24.
# DarkPlaces is the visible substrate; legacy Sierpinski/GEM/lore chrome must
# stay out of the frame until Quake-native ward composition lands.
ExecStart=
ExecStart=%h/.cache/hapax/source-activation/worktree/.venv/bin/python -m agents.studio_compositor --no-record --no-hls
Environment=CUDA_VISIBLE_DEVICES=1
Environment=HAPAX_3D_COMPOSITOR=0
Environment=HAPAX_V4L2_BRIDGE_ENABLED=1
Environment=HAPAX_DARKPLACES_BACKGROUND_ONLY=1
Environment=HAPAX_DARKPLACES_V4L2_DEVICE=/dev/video52
Environment=HAPAX_COMPOSITOR_STARTUP_EGRESS_GRACE_S=90
Environment=HAPAX_COMPOSITOR_DISABLE_V4L2_OUTPUT=0
Environment=HAPAX_COMPOSITOR_DISABLE_INLINE_FX=1
Environment=HAPAX_FX_AUTONOMOUS_MUTATIONS=0
Environment=HAPAX_SIERPINSKI_BASE_OVERLAY_ENABLED=0
Environment=HAPAX_SIERPINSKI_LAYOUT_SOURCE_ENABLED=0
Environment=HAPAX_PRE_FX_LAYOUT_DRAW_ENABLED=0
Environment=HAPAX_PRE_FX_LAYOUT_BACKGROUND_COMPOSITE_ENABLED=0
Environment=HAPAX_OVERLAY_ZONE_MANAGER_DRAW_ENABLED=0
Environment=HAPAX_GEM_SUBSTRATE_ENABLED=0
Environment=HAPAX_GEAL_ENABLED=0
Environment=HAPAX_HERO_SMALL_OVERLAY_ENABLED=0
Environment=HAPAX_DIRECTOR_SEGMENT_RUNNER_DISABLED=1
Environment=HAPAX_LORE_CHRONICLE_TICKER_ENABLED=0
Environment=HAPAX_LORE_PROGRAMME_STATE_ENABLED=0
Environment=HAPAX_LORE_RESEARCH_MARKER_CONTEXT_ENABLED=0
Environment=HAPAX_LORE_PRECEDENT_TICKER_ENABLED=0
Environment=HAPAX_LORE_PROGRAMME_HISTORY_ENABLED=0
Environment=HAPAX_LORE_RESEARCH_INSTRUMENT_DASHBOARD_ENABLED=0
Environment=HAPAX_LORE_RESEARCH_POSTER_CONSTRUCTIVIST_ENABLED=0
Environment=HAPAX_LORE_RESEARCH_POSTER_TUFTE_ENABLED=0
Environment=HAPAX_LORE_RESEARCH_POSTER_ASCII_ENABLED=0
Environment=HAPAX_V4L2_VIDEO42_WIDTH=640
Environment=HAPAX_V4L2_VIDEO42_HEIGHT=480
Environment=HAPAX_V4L2_VIDEO42_PIXEL_FORMAT=BGR4
Environment=HAPAX_V4L2_VIDEO42_FPS=30
Environment=HAPAX_V4L2_BRIDGE_WIDTH=640
Environment=HAPAX_V4L2_BRIDGE_HEIGHT=480
Environment=HAPAX_V4L2_BRIDGE_RAW_FORMAT=BGRx
Environment=HAPAX_V4L2_BRIDGE_PIXEL_FORMAT=BGR4
Environment=HAPAX_V4L2_BRIDGE_FPS=30
