# systemd/system/sshd.service.d/oom-protect.conf
[Service]
# Do not grant OOM immunity to SSH session children: oom_score_adj is inherited
# across fork/exec. Keep the daemon running if one session is selected instead.
OOMScoreAdjust=0
OOMPolicy=continue
