# systemd/system/sshd.service.d/oom-protect.conf
[Service]
# Keep SSH as the strongest remote recovery channel under host memory pressure.
OOMScoreAdjust=-1000
