# systemd/overrides/studio-compositor.service.d/v4l2-bridge.conf
[Service]
Environment=HAPAX_V4L2_BRIDGE_ENABLED=1
