# systemd/overrides/studio-compositor.service.d/v4l2-bridge.conf
[Service]
# Source rationale mirror for the deployable
# systemd/units/studio-compositor.service.d/v4l2-bridge.conf drop-in.
# Bridge fallback is opt-in while direct 3D V4L2 owns /dev/video42.
Environment=HAPAX_V4L2_BRIDGE_ENABLED=0
