# systemd/overrides/rag-ingest.service.d/gpu-isolate.conf
[Service]
# GPU isolation: docling's DocumentConverter loads a PyTorch layout/OCR model
# onto CUDA by default. Background RAG ingestion must never compete with
# TabbyAPI (EXL3 inference) or daimonion (STT) for VRAM. CPU is plenty for
# an agent that nobody is waiting on interactively.
# Same pattern as ollama.service — see workspace CLAUDE.md § Shared Infrastructure.
Environment=CUDA_VISIBLE_DEVICES=""

# Watch-only mode: disables the bulk_ingest() full-tree rescan that fires
# on every process start. The inotify watcher (ingest.py::watch) handles
# new/modified files as they land — drip-drip-drip instead of firehose.
# Bulk rescan can be performed manually with:
#   .venv-ingest/bin/python -m agents.ingest --bulk-only
ExecStart=
ExecStart=/home/hapax/projects/hapax-council/.venv-ingest/bin/python -m agents.ingest --watch-only
