# systemd/units/pipewire.service.d/cpu-affinity.conf
# Pin pipewire to logical cores {6, 7, 14, 15} on AMD Ryzen 7 7700X.
# See systemd/overrides/audio-stability/README.md for rationale.
[Service]
CPUAffinity=6 7 14 15
