# systemd/units/pipewire.service.d/oom-protect.conf
[Service]
# PipeWire carries NoNewPrivileges=yes, so sudo cannot cross the root boundary
# from ExecStartPost. Start valid; the root timer applies the protected score.
OOMScoreAdjust=100
MemoryLow=512M
MemoryMin=256M
