# systemd/units/pipewire.service.d/oom-protect.conf
[Service]
# Audio data-loop recovery must not depend on inherited user-manager OOM state.
OOMScoreAdjust=-900
MemoryLow=512M
MemoryMin=256M
