# systemd/units/pipewire.service.d/oom-protect.conf
[Service]
# User managers cannot lower child OOM scores. Start at the manager's valid
# score, then synchronously ask the root enforcer to apply the protected score.
OOMScoreAdjust=100
ExecStartPost=-/usr/local/bin/hapax-oom-score-trigger %n
MemoryLow=512M
MemoryMin=256M
