# systemd/units/youtube-player.service.d/slot-count.conf
[Service]
Environment=YOUTUBE_PLAYER_SLOTS=1
