# systemd/system/systemd-timesyncd.service.d/oom-protect.conf
[Service]
# Keep operator recovery channels alive during host memory pressure.
OOMScoreAdjust=-800
