# systemd/units/cache-cleanup.service
[Unit]
Description=Weekly cache cleanup — Docker, uv, pacman, stale worktree venvs

[Service]
Type=oneshot
ExecStart=%h/.local/bin/hapax-cache-cleanup
Environment=PATH=/home/hapax/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=/home/hapax
