# systemd/units/codex-claim-audit.service
[Unit]
Description=Codex claim audit (phantom + coherence + PR-reality check)
After=network-online.target
# See: scripts/codex-claim-audit --version

[Service]
Type=oneshot
WorkingDirectory=%h/.cache/hapax/source-activation/worktree
# Canonical-root: ExecStart uses the source-activation worktree, not the
# mutable dev tree. The activation worktree is a governed, reviewed checkout.
ExecStart=%h/.cache/hapax/source-activation/worktree/scripts/codex-claim-audit --release --stale-hours=6
Environment=GH_REPO=hapax-systems/hapax-council
Environment=NTFY_URL=https://ntfy.sh
Environment=NTFY_TOPIC=hapax-ops
# Reconcile cleared state to appendix (class #4 fix).
Environment=HAPAX_RECONCILE_HOST=hapax-appendix
