# systemd/units/codex-claim-audit.service
[Unit]
Description=Codex claim audit (phantom + coherence + PR-reality check)
# v2 makes the run network-dependent (gh pr view, ssh reconcile). Wants= pulls
# network-online.target in for a timer-launched oneshot so After= actually
# orders against it (After= alone does not start the target).
Wants=network-online.target
After=network-online.target
# See: scripts/codex-claim-audit --version

[Service]
Type=oneshot
WorkingDirectory=%h/.cache/hapax/source-activation/worktree
# Canonical-root: ExecStart uses the source-activation worktree, not the
# mutable dev tree. The activation worktree is a governed, reviewed checkout.
ExecStart=%h/.cache/hapax/source-activation/worktree/scripts/codex-claim-audit --release --stale-hours=6
Environment=GH_REPO=hapax-systems/hapax-council
Environment=NTFY_URL=https://ntfy.sh
Environment=NTFY_TOPIC=hapax-ops
# Reconcile cleared state to the PEER host (class #4 fix). This unit runs on
# podium (the single authoritative auditor), so the peer is appendix. The audit
# is vault-note-driven (the shared SSOT), so it reconciles the appendix replica
# for every active note — appendix-origin claims included. If this unit is ever
# also installed on appendix, point the reconcile host at hapax-podium there
# (never at self): a self-target is a no-op and would replicate nothing.
Environment=HAPAX_RECONCILE_HOST=hapax-appendix
