# systemd/units/gdrive-sync.timer
[Unit]
Description=gdrive sync (daily)

[Timer]
# Daily background reconciliation only. Fast file movement now uses direct
# KDE Connect/Tailscale/SSH paths instead of polling Drive as hot storage.
OnCalendar=*-*-* 18:30:00
RandomizedDelaySec=30min

[Install]
WantedBy=timers.target
