# systemd/units/hapax-voice-witness-watchdog.timer
[Unit]
Description=Run the voice witness watchdog every minute

[Timer]
OnBootSec=2min
OnUnitActiveSec=60s
AccuracySec=10s
Unit=hapax-voice-witness-watchdog.service

[Install]
WantedBy=timers.target
