# systemd/units/hapax-coord.service
[Unit]
Description=Hapax SBCL+CLOG coordination dashboard (source-only)
Documentation=file://%h/projects/hapax-coord/README.md
ConditionPathExists=%h/projects/hapax-coord/scripts/run-dev.sh

[Service]
Type=simple
WorkingDirectory=%h/projects/hapax-coord
Environment=HAPAX_COORD_COUNCIL_ROOT=%h/.cache/hapax/source-activation/worktree
ExecStart=%h/projects/hapax-coord/scripts/run-dev.sh --daemon --host 127.0.0.1 --port 8765
Restart=on-failure
RestartSec=5s
NoNewPrivileges=true

[Install]
WantedBy=default.target
