# systemd/units/hapax-broadcast-audio-health.service
[Unit]
Description=Hapax broadcast audio safety truth producer
Documentation=file:///home/hapax/projects/hapax-council/docs/superpowers/specs/2026-04-28-broadcast-audio-safety-ssot-design.md
After=pipewire.service pipewire-pulse.service wireplumber.service hapax-audio-safety.service hapax-private-broadcast-leak-guard.service
Wants=pipewire.service pipewire-pulse.service wireplumber.service hapax-audio-safety.service hapax-private-broadcast-leak-guard.service
OnFailure=notify-failure@%n.service

[Service]
Type=oneshot
WorkingDirectory=/home/hapax/projects/hapax-council
Environment=PATH=/home/hapax/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=/home/hapax
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/hapax/.local/bin/uv run --directory /home/hapax/projects/hapax-council python -m agents.broadcast_audio_health --loudness-duration 10
TimeoutStartSec=75s
MemoryMax=256M
CPUQuota=25%
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hapax-broadcast-audio-health

[Install]
WantedBy=default.target
