# systemd/units/hapax-broadcast-orchestrator.service
[Unit]
Description=Hapax YouTube Broadcast Boundary Orchestrator (CRITICAL — ytb-007)
Documentation=https://github.com/ryanklee/hapax-council/blob/main/agents/broadcast_orchestrator/__init__.py
After=network-online.target hapax-secrets.service hapax-private-broadcast-leak-guard.service
Wants=network-online.target hapax-private-broadcast-leak-guard.service
PartOf=hapax.target

[Service]
Type=notify
WorkingDirectory=%h/projects/hapax-council
EnvironmentFile=-/run/user/1000/hapax-secrets.env
Environment=PATH=%h/projects/hapax-council/.venv/bin:/usr/bin
Environment=HAPAX_BROADCAST_ORCHESTRATOR_ENABLED=0
Environment=HAPAX_BROADCAST_ROTATION_S=39600
Environment=HAPAX_BROADCAST_TICK_S=300
Environment=HAPAX_BROADCAST_PRIVACY_STATUS=public
Environment=HAPAX_BROADCAST_METRICS_PORT=9488
Environment=HAPAX_BROADCAST_STREAM_ID=
ExecStart=%h/projects/hapax-council/.venv/bin/python -m agents.broadcast_orchestrator
Restart=on-failure
RestartSec=30
WatchdogSec=600
MemoryMax=256M

[Install]
WantedBy=hapax.target
