# systemd/units/hapax-post-merge-deploy.service
[Unit]
Description=Auto-invoke hapax-post-merge-deploy when canonical main ref advances
Documentation=file://%h/projects/hapax-council/docs/superpowers/specs/2026-04-13-claude-md-excellence-design.md
After=network-online.target
OnFailure=notify-failure@%n.service

# Runaway-fire guard. The .path unit's PathChanged trigger only fires on
# refs/heads/main mtime changes, which the deploy script itself cannot
# cause (script writes to ~/.config/systemd/user/, ~/.config/pipewire/,
# ~/.local/bin/, and ~/.cache/hapax/post-merge-traces/ — none in .git).
# Belt-and-braces: cap to 3 invocations per 60s window so a runaway loop
# (caused by a future change) gets contained quickly and surfaced via
# OnFailure=notify-failure@%n.service. NOTE: StartLimitIntervalSec /
# StartLimitBurst belong on the [Unit] section, not [Service] — the
# 2026-05-02 waybar hardening bug pinned this exact landing site.
StartLimitIntervalSec=60
StartLimitBurst=3

# Companion to hapax-post-merge-deploy.path. Resolves the latest local main
# SHA at run time and hands it to scripts/hapax-post-merge-deploy, which
# classifies the diff and deploys the systemd-units / pipewire-confs /
# scripts side that rebuild-services.timer doesn't cover.

[Service]
Type=oneshot

WorkingDirectory=/home/hapax/projects/hapax-council
Environment=HOME=%h
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin

# Resolve the deploy target at run time. The .path unit fires when
# rebuild-service.sh advances refs/heads/main to origin/main (the canonical
# interactive worktree sits on a feature branch, so nothing else advances local
# main — the frozen + packed-only ref was the dead-trigger root cause). Resolve
# `origin/main` directly (not local `main`) so the deploy always targets the
# true merged HEAD even if the local ref lags. When a last-deployed-sha receipt
# exists and is an ancestor, deploy the CUMULATIVE --since range so a multi-merge
# jump installs every intermediate merge's units/confs/scripts, not just the tip.
ExecStart=/bin/bash -c 'repo=%h/projects/hapax-council; sha=$(git -C "$repo" rev-parse origin/main 2>/dev/null); [ -n "$sha" ] || exit 0; since=$(cat %h/.cache/hapax/post-merge-traces/last-deployed-sha 2>/dev/null || true); if [ -n "$since" ] && [ "$since" != "$sha" ] && git -C "$repo" merge-base --is-ancestor "$since" "$sha" 2>/dev/null; then exec %h/.local/bin/hapax-post-merge-deploy --since "$since" "$sha"; else exec %h/.local/bin/hapax-post-merge-deploy "$sha"; fi'

StandardOutput=journal
StandardError=journal
SyslogIdentifier=hapax-post-merge-deploy

# Guard rails — the deploy script itself manages systemctl --user restarts
# and writes small text files; no GPU/big-memory resident state. Keep bounded
# but do not let the deploy rail fail solely because a 128G host inherited an
# old 1G hard ceiling.
MemoryHigh=1G
MemoryMax=2G
TimeoutStartSec=180
