# systemd/units/hapax-hooks-doctor.timer
[Unit]
Description=Periodic Hapax gate-drift sweep schedule (reform FM-6)
# Hapax-Auto-Enable: true

[Timer]
OnBootSec=10min
OnUnitActiveSec=6h
RandomizedDelaySec=10min
Persistent=true
Unit=hapax-hooks-doctor.service

[Install]
WantedBy=timers.target
