# systemd/units/hapax-dataset-card-generator.service
[Unit]
Description=Hapax dataset card generator
After=hapax-secrets.service
ConditionPathExists=%h/projects/hapax-council/config/research-corpus-export-ledger.yaml

[Service]
Type=oneshot
WorkingDirectory=%h/.cache/hapax/source-activation/worktree
Environment=HOME=%h
Environment=PYTHONPATH=%h/.cache/hapax/source-activation/worktree
ExecStart=%h/.cache/hapax/source-activation/worktree/.venv/bin/python -m agents.dataset_card_generator --output %h/hapax-state/research/dataset-cards.md

[Install]
WantedBy=default.target
