# systemd/units/hapax-equipment-state.service
[Unit]
Description=Hapax equipment state writer
Documentation=config/equipment/
ConditionPathExists=%h/.cache/hapax/source-activation/worktree/scripts/equipment-state-writer

[Service]
Type=oneshot
ExecStart=%h/.cache/hapax/source-activation/worktree/scripts/equipment-state-writer
Environment=PYTHONPATH=%h/.cache/hapax/source-activation/worktree
TimeoutSec=30
