# systemd/units/hapax-visual-pool-snapshot-harvester.service
[Unit]
Description=Hapax visual pool snapshot harvester
After=studio-compositor.service
ConditionPathIsDirectory=/dev/shm/hapax-compositor

[Service]
Type=oneshot
WorkingDirectory=%h/projects/hapax-council
ExecStart=%h/.local/bin/uv --directory %h/projects/hapax-council run python scripts/visual-pool-snapshot-harvester.py
Environment=PATH=%h/.local/bin:%h/.cargo/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=%h
MemoryMax=256M
CPUQuota=25%
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hapax-visual-pool-snapshot-harvester
