# systemd/units/hapax-audio-perception.service
[Unit]
Description=Hapax audio perception daemon — CPU-only scene classification
After=pipewire.service pipewire-pulse.service hapax-secrets.service
Wants=pipewire-pulse.service
OnFailure=notify-failure@%n.service
ConditionPathExists=%h/.cache/hapax/rebuild/worktree/pyproject.toml

[Service]
Type=simple
WorkingDirectory=%h/.cache/hapax/rebuild/worktree
ExecStart=%h/.local/bin/uv run --extra audio python -m agents.audio_perception
Restart=always
RestartSec=10
EnvironmentFile=/run/user/1000/hapax-secrets.env
Environment=PATH=/home/hapax/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=/home/hapax
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=CUDA_VISIBLE_DEVICES=
Environment=OMP_NUM_THREADS=2
Environment=MKL_NUM_THREADS=2
Environment=OPENBLAS_NUM_THREADS=2
MemoryMax=4G
MemoryHigh=3G
CPUQuota=200%
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hapax-audio-perception

[Install]
WantedBy=default.target
