# systemd/units/hapax-chronicle-high-salience-public-event-producer.service
[Unit]
Description=Hapax chronicle high-salience public-event producer
After=hapax-secrets.service
PartOf=hapax.target

[Service]
Type=simple
WorkingDirectory=%h/.cache/hapax/source-activation/worktree
EnvironmentFile=-/run/user/1000/hapax-secrets.env
Environment=PATH=%h/.cache/hapax/source-activation/worktree/.venv/bin:%h/.local/bin:/usr/local/bin:/usr/bin
Environment=PYTHONPATH=%h/.cache/hapax/source-activation/worktree
Environment=LOG_LEVEL=INFO
Environment=HAPAX_CHRONICLE_EVENTS_PATH=/dev/shm/hapax-chronicle/events.jsonl
Environment=HAPAX_RESEARCH_VEHICLE_PUBLIC_EVENT_PATH=/dev/shm/hapax-public-events/events.jsonl
ExecStart=%h/.cache/hapax/source-activation/worktree/.venv/bin/python -m agents.chronicle_high_salience_public_event_producer
Restart=on-failure
RestartSec=15
MemoryMax=128M

[Install]
WantedBy=default.target
