# systemd/units/hapax-oom-policy-audit.timer
# Hapax-Installer-Owner: scripts/install-p0-oom-containment
[Unit]
Description=Run Hapax OOM containment policy audit
Documentation=https://github.com/hapax-systems/hapax-council/blob/main/systemd/README.md

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Unit=hapax-oom-policy-audit.service

[Install]
WantedBy=timers.target
