# systemd/units/hapax-self-federate-rss.service
[Unit]
Description=Hapax weblog RSS validator — weekly self-federate feed health check
After=hapax-secrets.service

[Service]
Type=oneshot
ExecStart=%h/.cache/hapax/source-activation/worktree/.venv/bin/python -m agents.self_federate.rss_validator
WorkingDirectory=%h/.cache/hapax/source-activation/worktree
Environment=HOME=%h
Environment=PYTHONPATH=%h/.cache/hapax/source-activation/worktree

[Install]
WantedBy=default.target
