# systemd/units/hapax-live-surface-guard.service
[Unit]
Description=Live Surface Guard — OBS/HLS/public egress observability and bounded remediation receipts
After=studio-compositor.service hapax-hls-no-cache.service
Wants=hapax-hls-no-cache.service
ConditionPathExists=%h/.cache/hapax/source-activation/worktree/agents/live_surface_guard

[Service]
Type=simple
WorkingDirectory=%h/.cache/hapax/source-activation/worktree
ExecStartPre=%h/.cache/hapax/source-activation/worktree/scripts/hapax-compositor-runtime-source-check --require-file agents/live_surface_guard/__main__.py --require-file shared/live_surface_truth.py
ExecStart=%h/.cache/hapax/source-activation/worktree/.venv/bin/python -m agents.live_surface_guard --require-hls --require-obs-decoder --textfile-path %h/.local/share/node_exporter/textfile_collector/hapax-live-surface-guard.prom
Restart=on-failure
RestartSec=5s
EnvironmentFile=-%t/hapax-secrets.env
Environment=PYTHONPATH=%h/.cache/hapax/source-activation/worktree
Environment=HAPAX_COMPOSITOR_SOURCE_ROOT=%h/.cache/hapax/source-activation/worktree

[Install]
WantedBy=default.target
