# systemd/units/hapax-cc-task-offer-ready.service
[Unit]
Description=Hapax CC task ready-to-offered promotion reconciler
After=hapax-secrets.service

[Service]
Type=oneshot
ExecStart=%h/projects/hapax-council/scripts/cc-task-offer-ready --reconcile
Environment=HOME=%h
TimeoutStartSec=60
