# systemd/units/hapax-datacite-graph-publish.service
[Unit]
Description=Hapax DataCite citation-graph DOI publisher — Phase 2B chained mint after daily mirror refresh
After=hapax-secrets.service hapax-datacite-mirror.service
Wants=hapax-datacite-mirror.service

[Service]
Type=oneshot
ExecStart=%h/.cache/hapax/source-activation/worktree/.venv/bin/python -m agents.publication_bus.self_citation_graph_doi --commit
WorkingDirectory=%h/.cache/hapax/source-activation/worktree
Environment=HOME=%h
Environment=PYTHONPATH=%h/.cache/hapax/source-activation/worktree
EnvironmentFile=-%h/.cache/hapax/secrets.env

[Install]
WantedBy=default.target
