# systemd/units/hapax-dmn.service
[Unit]
Description=Hapax DMN — always-on cognitive substrate
After=hapax-secrets.service ollama.service
Wants=ollama.service
PartOf=hapax-visual-stack.target

[Service]
Type=simple
WorkingDirectory=%h/.cache/hapax/source-activation/worktree
ExecStart=%h/.cache/hapax/source-activation/worktree/.venv/bin/python -m agents.dmn
Restart=on-failure
RestartSec=10
EnvironmentFile=/run/user/1000/hapax-secrets.env
Environment=HAPAX_SOURCE_ACTIVATION_WORKTREE=%h/.cache/hapax/source-activation/worktree
Environment=PYTHONPATH=%h/.cache/hapax/source-activation/worktree
Environment=PATH=%h/.cache/hapax/source-activation/worktree/.venv/bin:%h/.local/bin:/usr/local/bin:/usr/bin
MemoryHigh=2G
MemoryMax=4G

[Install]
WantedBy=default.target
