# systemd/units/hapax-broadcast-orchestrator.service
[Unit]
Description=Hapax YouTube Broadcast Boundary Orchestrator (CRITICAL — ytb-007)
Documentation=https://github.com/ryanklee/hapax-council/blob/main/agents/broadcast_orchestrator/__init__.py
After=network-online.target hapax-secrets.service hapax-private-broadcast-leak-guard.service
Wants=network-online.target hapax-private-broadcast-leak-guard.service
PartOf=hapax.target
OnFailure=notify-failure@%n.service

[Service]
Type=notify
WorkingDirectory=%h/.cache/hapax/source-activation/worktree
EnvironmentFile=-/run/user/1000/hapax-secrets.env
Environment=PATH=%h/.cache/hapax/source-activation/worktree/.venv/bin:%h/.local/bin:/usr/bin
Environment=PYTHONPATH=%h/.cache/hapax/source-activation/worktree
Environment=HAPAX_BROADCAST_ORCHESTRATOR_ENABLED=0
Environment=HAPAX_BROADCAST_ROTATION_S=39600
Environment=HAPAX_BROADCAST_TICK_S=300
Environment=HAPAX_BROADCAST_PRIVACY_STATUS=public
Environment=HAPAX_BROADCAST_METRICS_PORT=9488
Environment=HAPAX_BROADCAST_STREAM_ID=
ExecStartPre=%h/.cache/hapax/source-activation/worktree/scripts/hapax-compositor-runtime-source-check --require-file agents/broadcast_orchestrator/__main__.py
ExecStart=%h/.cache/hapax/source-activation/worktree/.venv/bin/python -m agents.broadcast_orchestrator
Restart=on-failure
RestartSec=30
WatchdogSec=600
MemoryMax=256M

[Install]
WantedBy=hapax.target
