# systemd/units/hapax-datacite-snapshot.service
[Unit]
Description=Hapax DataCite citation-graph snapshot — nightly graph network around DOI/SWHID/ORCID nodes
After=hapax-secrets.service

[Service]
Type=oneshot
ExecStart=%h/.cache/hapax/source-activation/worktree/.venv/bin/python -m agents.attribution.datacite_graphql_snapshot
WorkingDirectory=%h/.cache/hapax/source-activation/worktree
Environment=HOME=%h
Environment=PYTHONPATH=%h/.cache/hapax/source-activation/worktree
EnvironmentFile=-%h/.cache/hapax/secrets.env
EnvironmentFile=-%h/.config/hapax/datacite-snapshot.env

[Install]
WantedBy=default.target
