# systemd/units/hapax-backup-gdrive-critical.timer
[Unit]
Description=Daily GDrive critical offsite backup

[Timer]
OnCalendar=*-*-* 04:35:00
RandomizedDelaySec=45m
Persistent=true

[Install]
WantedBy=timers.target
