# systemd/units/hapax-information-density.service
[Unit]
Description=Hapax Information Density Field — universal perception-to-automation primitive
After=hapax-secrets.service pipewire.service
Wants=hapax-secrets.service

[Service]
Type=simple
ExecStartPre=/bin/sleep 5
ExecStart=%h/.cache/hapax/source-activation/worktree/.venv/bin/python -m agents.information_density_daemon
WorkingDirectory=%h/.cache/hapax/source-activation/worktree
Environment=PYTHONPATH=%h/.cache/hapax/source-activation/worktree
Restart=on-failure
RestartSec=10
MemoryMax=256M
WatchdogSec=120

[Install]
WantedBy=default.target
