# systemd/units/hapax-mastodon-post.service
[Unit]
Description=Hapax Mastodon poster (canonical public-event federation)
Documentation=https://github.com/ryanklee/council/blob/main/agents/cross_surface/mastodon_post.py
After=network-online.target hapax-secrets.service
Wants=network-online.target
PartOf=hapax.target

[Service]
Type=simple
WorkingDirectory=%h/.cache/hapax/source-activation/worktree
EnvironmentFile=-/run/user/1000/hapax-secrets.env
Environment=PATH=%h/.cache/hapax/source-activation/worktree/.venv/bin:%h/.local/bin:/usr/local/bin:/usr/bin
Environment=PYTHONPATH=%h/.cache/hapax/source-activation/worktree
Environment=HAPAX_MASTODON_METRICS_PORT=9502
Environment=HAPAX_MASTODON_TICK_S=30
Environment=HAPAX_LOG_LEVEL=INFO
# Tails canonical ResearchVehiclePublicEvent JSONL records from:
#   HAPAX_RESEARCH_VEHICLE_PUBLIC_EVENT_PATH
#   default: /dev/shm/hapax-public-events/events.jsonl
# Cursor/idempotency:
#   HAPAX_MASTODON_CURSOR
#   HAPAX_MASTODON_IDEMPOTENCY_PATH
# Metric:
#   hapax_broadcast_mastodon_posts_total{result}
# Operator must export from hapax-secrets:
#   HAPAX_MASTODON_INSTANCE_URL    (e.g. https://mastodon.social)
#   HAPAX_MASTODON_ACCESS_TOKEN    (from <instance>/settings/applications, scope write:statuses)
# Without either, daemon idles + logs "no_credentials" per event.
ExecStart=%h/.cache/hapax/source-activation/worktree/.venv/bin/python -m agents.cross_surface.mastodon_post
Restart=on-failure
RestartSec=15
MemoryMax=128M

[Install]
WantedBy=default.target
