# systemd/units/hapax-lane-idle-watchdog.service
[Unit]
Description=Hapax Lane Idle Watchdog — pane-level idle detection + auto-redispatch
After=network.target

[Service]
Type=oneshot
ExecStart=%h/.cache/hapax/source-activation/worktree/scripts/hapax-lane-idle-watchdog
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=%h
Environment=HAPAX_LOCAL_DEV_MAINTENANCE_MODE=appendix-only
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hapax-lane-idle-watchdog
TimeoutStartSec=90

[Install]
WantedBy=default.target
