# systemd/units/hapax-oom-policy-audit.service
[Unit]
Description=Audit Hapax OOM containment policy drift
Documentation=https://github.com/hapax-systems/hapax-council/blob/main/systemd/README.md
OnFailure=notify-failure@%n.service
StartLimitIntervalSec=0

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/hapax-oom-policy-audit --json
SyslogIdentifier=hapax-oom-policy-audit
