# systemd/units/hapax-opus-route-authority-receipt.service
[Unit]
Description=Maintain a fresh operator-signed opus route-authority receipt (standing OQ-5)
Documentation=file://%h/.cache/hapax/source-activation/worktree/docs/routing-ontology-reference.md
OnFailure=notify-failure@%n.service

[Service]
Type=oneshot
# OQ-5 made standing: re-sign the operator's opus model-entitlement so the live
# dispatch read-path discovers a fresh route-authority receipt in the default
# receipt dir (~/.cache/hapax/platform-capability-receipts/route-authority/) and
# un-degrades claude.headless.opus to LAUNCH — no --policy-rollback, no
# dispatch-from-a-stale-worktree hack. --ensure-fresh converges on a single
# stable-id receipt and only re-signs once the 24h window drops below 8h
# remaining (the paired .timer fires every 6h, so it never lapses).
#
# Source worktree = the source-activation ACTIVE deploy symlink, NOT the primary
# `projects/hapax-council` worktree. The primary lane parks on feature branches
# (e.g. alpha/screwm-*) where this minter does not exist, which silently broke
# upkeep (status=2 "No such file") until the receipt lapsed. The symlink always
# resolves to the activated origin/main release (script + synced .venv present)
# and its PATH is stable across deploys — only its target advances.
ExecStart=%h/.cache/hapax/source-activation/worktree/.venv/bin/python %h/.cache/hapax/source-activation/worktree/scripts/hapax-mint-route-authority-receipt --ensure-fresh --receipt-type opus_model_entitlement --route-id claude.headless.opus --stale-after 24h --refresh-within 8h --alert-within 2h --alert-after-failures 3
WorkingDirectory=%h/.cache/hapax/source-activation/worktree
Environment=HOME=%h
Environment=PYTHONPATH=%h/.cache/hapax/source-activation/worktree
# A failed re-mint is suppressed by the notify-failure coalescer (timer-backed),
# so the minter self-escalates an un-suppressed ntfy via these on N consecutive
# failures or within --alert-within of expiry.
Environment=NTFY_URL=https://ntfy.sh
Environment=NTFY_TOPIC=hapax-ops
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hapax-opus-route-authority-receipt
MemoryMax=256M
TimeoutStartSec=60
# Killswitch: `systemctl --user mask hapax-opus-route-authority-receipt.timer`
# stops upkeep; the live receipt then lapses after --stale-after and opus HOLDs.
