# systemd/units/hapax-lane-supervisor.timer
# Hapax-Auto-Enable: true
# The recurring surface is safe to auto-enable because it only projects lane
# health candidates and cannot launch, signal, or detach ownership.
[Unit]
Description=Project Hapax lane health every 60 seconds

[Timer]
OnBootSec=60
OnUnitActiveSec=60
AccuracySec=10

[Install]
WantedBy=timers.target
