# systemd/units/hapax-research-registry-producer.service
[Unit]
Description=Hapax research-registry producer — journals research artefacts to ~/hapax-state/research/registry.jsonl
After=hapax-secrets.service

[Service]
Type=oneshot
ExecStart=%h/projects/hapax-council/.venv/bin/python %h/projects/hapax-council/scripts/research-registry-emit.py
WorkingDirectory=%h/projects/hapax-council
Environment=HOME=%h
Environment=PYTHONPATH=%h/projects/hapax-council

[Install]
WantedBy=default.target
