# systemd/units/hapax-audio-self-perception.service
[Unit]
Description=Hapax audio self-perception loop — AVSDLC-002 broadcast-master spectral tap
After=pipewire.service pipewire-pulse.service
Requires=pipewire-pulse.service

[Service]
Type=simple
ExecStart=%h/.local/share/hapax-council/.venv/bin/python -m agents.audio_self_perception
WorkingDirectory=%h/.local/share/hapax-council
Restart=on-failure
RestartSec=10
WatchdogSec=60s
MemoryMax=256M

[Install]
WantedBy=default.target
