# systemd/units/hapax-novelty-shift-emitter.service
[Unit]
Description=Hapax novelty-shift impingement emitter (long-running loop)
Documentation=https://github.com/ryanklee/hapax-council/blob/main/agents/novelty_emitter/_emitter.py
After=hapax-daimonion.service hapax-secrets.service
StartLimitIntervalSec=0

[Service]
Type=simple
WorkingDirectory=%h/.cache/hapax/source-activation/worktree
ExecStart=/usr/bin/uv run python -m agents.novelty_emitter --loop --interval 1
Restart=always
RestartSec=2
TimeoutStartSec=10

[Install]
WantedBy=default.target
