# systemd/units/hapax-audio-routing-check.service
[Unit]
Description=Hapax audio routing invariant check
After=pipewire.service

[Service]
Type=oneshot
ExecStart=%h/.local/bin/hapax-audio-routing-check
StandardOutput=journal
StandardError=journal
# Exit 1 = invariant violation, exit 2 = PipeWire not running.
# Both are genuine failures that must surface to systemd so
# 'systemctl --user status' and alerting see them.
