# systemd/units/hapax-orcid-verifier.service
[Unit]
Description=Hapax ORCID auto-update verifier — daily check that minted concept-DOIs flow to operator's ORCID record
After=hapax-secrets.service

[Service]
Type=oneshot
ExecStart=%h/.cache/hapax/source-activation/worktree/.venv/bin/python -m agents.publication_bus.orcid_verifier
WorkingDirectory=%h/.cache/hapax/source-activation/worktree
Environment=HOME=%h
Environment=PYTHONPATH=%h/.cache/hapax/source-activation/worktree
EnvironmentFile=-%h/.cache/hapax/secrets.env

[Install]
WantedBy=default.target
