# systemd/units/hapax-cc-pr-autoqueue.timer
[Unit]
Description=Hapax governed PR auto-queue reconciler timer

[Timer]
OnBootSec=4min
OnUnitActiveSec=6min
RandomizedDelaySec=45s
AccuracySec=30s

[Install]
WantedBy=timers.target
