# systemd/units/hapax-screwm-camera-gamepad.service
[Unit]
Description=Hapax Screwm headless gamepad camera bridge
PartOf=hapax-darkplaces-v4l2.service
After=hapax-darkplaces-v4l2.service
ConditionPathExists=%h/.config/hapax/enable-darkplaces-runtime
ConditionPathExists=%h/.config/hapax/enable-screwm-camera-gamepad

[Service]
Type=simple
WorkingDirectory=%h/.cache/hapax/source-activation/worktree
ExecStartPre=%h/.cache/hapax/source-activation/worktree/scripts/hapax-compositor-runtime-source-check --require-file scripts/screwm-camera-gamepad.py
ExecStart=/usr/bin/bash -lc 'exec "$HOME/.cache/hapax/source-activation/worktree/scripts/screwm-camera-gamepad.py" --device /dev/input/js2'
Restart=on-failure
RestartSec=2

[Install]
WantedBy=hapax-darkplaces-v4l2.service
