# systemd/units/hapax-s4-rearm.timer
[Unit]
Description=Re-assert the S-4 gain ladder silently at boot + every 2 minutes
# No Requires=hapax-s4-rearm.service: a .timer activates its same-named .service
# on schedule by default. A Requires= here would also pull the service the moment
# the timer starts at boot, bypassing the OnBootSec=60s settle delay and exposing
# it to the cold-boot enumeration race (review finding, gemini-1).

[Timer]
# OnBootSec=60s lets PipeWire/WirePlumber + the S-4 USB-MIDI enumerate before
# the first silent re-assert. OnUnitActiveSec=120s keeps the ladder asserted
# across the boot session, so any daimonion/router restart that left the S-4
# unarmed is re-armed within ~2 min. Silent (--witness-mode none) ⇒ no tone to
# air; the recurrence also covers the cold-boot enumeration race (a first-fire
# midi_output_missing self-heals on the next tick).
OnBootSec=60s
OnUnitActiveSec=120s
AccuracySec=10s
# Persistent= is a no-op on monotonic timers (OnBootSec/OnUnitActiveSec, no
# OnCalendar) — kept only to match the hapax-audio-topology-assertion.timer
# precedent. OnUnitActiveSec is what actually carries the recurrence.
Persistent=true

[Install]
WantedBy=timers.target
