# systemd/units/hapax-cc-pr-merge-watcher.timer
[Unit]
Description=Hapax CC-task PR-merge auto-close watcher timer

[Timer]
OnBootSec=3min
OnUnitActiveSec=9min
RandomizedDelaySec=45s
AccuracySec=30s

[Install]
WantedBy=timers.target
