# systemd/units/hapax-oom-score-enforce.timer
[Unit]
# Hapax-Install-Scope: system
Description=Reapply live Hapax OOM scores
Documentation=file:/home/hapax/.local/lib/hapax-recovery/council/current/manifest.json

[Timer]
OnBootSec=30s
OnUnitActiveSec=30s
AccuracySec=1s
Unit=hapax-oom-score-enforce.service

[Install]
WantedBy=timers.target
