# systemd/units/hapax-oom-score-enforce.timer
[Unit]
# Hapax-Install-Scope: system
Description=Reapply live Hapax OOM scores
Documentation=file:/home/hapax/.cache/hapax/source-activation/worktree/systemd/README.md
ConditionPathExists=/usr/local/sbin/hapax-oom-score-enforce

[Timer]
OnBootSec=30s
OnUnitActiveSec=30s
Unit=hapax-oom-score-enforce.service

[Install]
WantedBy=timers.target
