# systemd/units/hapax-darkplaces-bridge.service
[Unit]
Description=Hapax DarkPlaces state bridge — copies live state into game directory
PartOf=hapax-darkplaces.service
After=hapax-darkplaces.service
ConditionPathExists=%h/.config/hapax/enable-darkplaces-runtime

[Service]
Type=simple
WorkingDirectory=%h/.cache/hapax/source-activation/worktree
ExecStartPre=%h/.cache/hapax/source-activation/worktree/scripts/hapax-compositor-runtime-source-check --require-file scripts/darkplaces-state-bridge.sh
ExecStart=/usr/bin/bash -lc 'exec "$HOME/.cache/hapax/source-activation/worktree/scripts/darkplaces-state-bridge.sh"'
Restart=always
RestartSec=2

[Install]
WantedBy=hapax-darkplaces.service
