# systemd/units/hapax-omg-weblog-composer.service
[Unit]
Description=Hapax omg.lol weblog composer — monthly draft skeleton emitter to ~/hapax-state/weblog-drafts/
After=hapax-secrets.service

[Service]
Type=oneshot
ExecStart=%h/.cache/hapax/source-activation/worktree/.venv/bin/python -m agents.omg_weblog_composer
WorkingDirectory=%h/.cache/hapax/source-activation/worktree
Environment=HOME=%h
Environment=PYTHONPATH=%h/.cache/hapax/source-activation/worktree

[Install]
WantedBy=default.target
