# systemd/units/hapax-lane-reaper.service
[Unit]
Description=Hapax lane recovery projection (no session, worktree, task, or claim effects)

[Service]
Type=oneshot
ExecStart=/home/hapax/projects/hapax-council/scripts/hapax-lane-reaper --threshold 30
Environment=HOME=/home/hapax
TimeoutStartSec=120
