# systemd/units/hapax-lane-reaper.service
[Unit]
Description=Hapax lane recovery projection (no session, worktree, task, or claim effects)
After=hapax-source-activate.service
ConditionFileIsExecutable=%h/.cache/hapax/source-activation/worktree/scripts/hapax-lane-reaper

[Service]
Type=oneshot
WorkingDirectory=%h
ExecStart=/usr/bin/bash -c 'set -euo pipefail; state="%h/.cache/hapax/source-activation"; IFS= read -r sha < "$state/last-success-sha"; case "$sha" in ""|*[!0-9a-f]*) echo "hapax-lane-reaper: invalid activation receipt" >&2; exit 78;; esac; [ "${#sha}" -eq 40 ] || { echo "hapax-lane-reaper: invalid activation receipt" >&2; exit 78; }; target="$(/usr/bin/readlink -f "$state/worktree")"; [ "${target##*/}" = "$sha" ] && [ "$(/usr/bin/git -C "$target" rev-parse HEAD)" = "$sha" ] && [ -z "$(/usr/bin/git -C "$target" status --porcelain --untracked-files=no -- scripts/hapax-lane-reaper shared)" ] || { echo "hapax-lane-reaper: activated release identity mismatch" >&2; exit 78; }; exec "$target/scripts/hapax-lane-reaper" --threshold 30'
Environment=HOME=%h
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=PYTHONDONTWRITEBYTECODE=1
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hapax-lane-reaper
TimeoutStartSec=120
NoNewPrivileges=true
CapabilityBoundingSet=
AmbientCapabilities=
ProtectSystem=strict
ProtectHome=read-only
ReadOnlyPaths=%t /tmp /var/tmp
PrivateNetwork=true
RestrictAddressFamilies=AF_UNIX
PrivateTmp=true
PrivateDevices=true
ProtectControlGroups=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectClock=true
RestrictNamespaces=true
RestrictSUIDSGID=true
LockPersonality=true
SystemCallArchitectures=native
SystemCallFilter=~@network-io kill pidfd_send_signal rt_sigqueueinfo rt_tgsigqueueinfo tkill tgkill
