# systemd/units/hapax-v4l2-bridge.service
[Unit]
Description=V4L2 Bridge — shmsrc→v4l2sink sidecar for /dev/video42
After=studio-compositor.service
BindsTo=studio-compositor.service
PartOf=studio-compositor.service

[Service]
Type=simple
ExecStart=%h/projects/hapax-council/scripts/hapax-v4l2-bridge
Restart=always
RestartSec=1s
WatchdogSec=30s
Environment=GST_DEBUG=2
Environment=HAPAX_V4L2_BRIDGE_DEVICE=/dev/video42
Environment=HAPAX_V4L2_BRIDGE_SOCKET=/dev/shm/hapax-compositor/v4l2-bridge.sock

[Install]
WantedBy=default.target
