# systemd/units/hapax-backup-gdrive-critical.service
[Unit]
Description=Hapax GDrive Critical Offsite Backup (restic)
After=network-online.target
Wants=network-online.target
OnFailure=notify-failure@%n.service
ConditionPathExists=/store/llm-data
ConditionPathExists=%h/.config/rclone/rclone.conf

[Service]
Type=oneshot
ExecStart=/home/hapax/projects/hapax-council/scripts/hapax-backup-gdrive-critical
WorkingDirectory=/home/hapax/projects/hapax-council
Environment=PATH=/home/hapax/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=/home/hapax
Environment=GNUPGHOME=/home/hapax/.gnupg
Environment=PASSWORD_STORE_DIR=/home/hapax/.password-store
Environment=HAPAX_GDRIVE_CRITICAL_RESTIC_PASSWORD_ENTRY=backblaze/restic-password
Environment=RESTIC_CACHE_DIR=/store/llm-data/restic-cache/gdrive-critical
MemoryHigh=1G
MemoryMax=2G
CPUQuota=25%
TimeoutStartSec=2h
SyslogIdentifier=hapax-backup-gdrive-critical
