# systemd/units/hapax-audio-health-inter-stage-corr.service
[Unit]
Description=Hapax audio health M4 — inter-stage envelope correlation monitor
After=pipewire.service pipewire-pulse.service
Requires=pipewire-pulse.service
PartOf=hapax-audio-health.target
ConditionPathExists=%h/.cache/hapax/rebuild/worktree/pyproject.toml

[Service]
Type=simple
NotifyAccess=all
WorkingDirectory=%h/.cache/hapax/rebuild/worktree
ExecStart=%h/.local/bin/uv run python -m agents.audio_health.m4_inter_stage_corr_daemon
Restart=on-failure
RestartSec=10
WatchdogSec=0
EnvironmentFile=-%h/.config/hapax/audio-health.env
MemoryMax=256M

[Install]
WantedBy=hapax-audio-health.target
