# systemd/units/hapax-cred-watch.service
[Unit]
Description=Hapax credential watch — snapshot pass entry names; emit operator-unblocker report (no values, ever)
After=hapax-secrets.service

[Service]
Type=oneshot
ExecStart=%h/.cache/hapax/source-activation/worktree/.venv/bin/python -m agents.hapax_cred_monitor --once
WorkingDirectory=%h/.cache/hapax/source-activation/worktree
Environment=HOME=%h
Environment=PYTHONPATH=%h/.cache/hapax/source-activation/worktree

[Install]
WantedBy=default.target
