# systemd/units/hapax-rtsp-loopback-watchdog.timer
[Unit]
Description=Poll RTSP loopback producer health

[Timer]
OnBootSec=60
OnUnitActiveSec=10
AccuracySec=1s

[Install]
WantedBy=timers.target
