# systemd/units/hapax-m8-stem-retention.timer
[Unit]
Description=Hapax M8 stem retention sweep — daily

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=10m

[Install]
WantedBy=timers.target
