# systemd/units/hapax-storage-pressure-check.service
[Unit]
Description=Hapax root Btrfs and filesystem storage pressure check
Documentation=file:%h/projects/hapax-council/scripts/hapax-storage-pressure-check

[Service]
Type=oneshot
ExecStart=%h/projects/hapax-council/scripts/hapax-storage-pressure-check
Environment=HOME=%h
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
TimeoutStartSec=60
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hapax-storage-pressure
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
