# systemd/units/hapax-option-c-pin-watchdog.service
[Unit]
Description=Hapax Option-C HDA pin-watchdog (Auditor A O3c)
Documentation=https://github.com/ryanklee/hapax-council/blob/main/docs/runbooks/audio-incidents.md#ryzen-hda-pin-glitch
After=pipewire.service wireplumber.service pipewire-pulse.service

[Service]
Type=oneshot
ExecStart=%h/.local/bin/hapax-option-c-repin
TimeoutStartSec=10
SuccessExitStatus=0

[Install]
WantedBy=default.target
