# systemd/units/hapax-dispatch-redemption.service
[Unit]
# Hapax-Install-Scope: system
#
# System-scoped on purpose (scout audit, cc-task-dispatch-redemption-runtime-governor-20260704):
# the fixed governor socket lives under system /run (/run/hapax/coord), which a
# user unit cannot own — user units get /run/user/<uid>, and a caller-writable
# /run/user path is exactly the same-UID surface the redemption design forbids.
# RuntimeDirectory= below makes the system manager provision a root-owned
# namespace with deliberate group/mode (no ad hoc mkdir), and remove it when
# the governor stops, so wrappers fail closed on connect instead of trusting a
# stale inode. The socket is root:hapax 0660: hapax can connect but cannot
# unlink/replace the socket path in the root-owned 0750 directory.
#
# Install/activate:
#   scripts/hapax-dispatch-redemption-service-install --install
#
# hapax-post-merge-deploy invokes that installer automatically when this unit
# changes, because this root/system-scoped service must not be copied into the
# user unit directory and must not merge dormant.
#
# Deployment root:
#   This unit is intentionally pinned to the canonical single-operator
#   `/home/hapax/.cache/hapax/source-activation/worktree` activation root.
#   Alternate users/hosts must regenerate this unit via the installer after
#   changing the allowed user, activation paths, and requester digest together.
#
# Witness owner/mode + live protocol afterwards (token-free receipt):
#   scripts/hapax-dispatch-redemption-authority --receipt
Description=Hapax dispatch launch redemption governor (fixed-socket mint/redeem authority)
Documentation=file:/home/hapax/.cache/hapax/source-activation/worktree/docs/methodology/emergency-bypass.md
ConditionPathExists=/home/hapax/.cache/hapax/source-activation/worktree/scripts/hapax-dispatch-redemption-authority

[Service]
Type=simple
User=root
Group=hapax
Environment=HAPAX_DISPATCH_REDEMPTION_ALLOWED_USER=hapax
Environment=HAPAX_DISPATCH_REDEMPTION_ALLOWED_REQUESTER=hapax-methodology-dispatch
Environment=HAPAX_DISPATCH_REDEMPTION_ALLOWED_REQUESTER_PATHS=/home/hapax/.cache/hapax/source-activation/worktree/scripts/hapax-methodology-dispatch
# Refresh this digest whenever scripts/hapax-methodology-dispatch changes;
# recheck command:
#   sha256sum /home/hapax/.cache/hapax/source-activation/worktree/scripts/hapax-methodology-dispatch
# tests/systemd/test_dispatch_redemption_unit.py enforces the repo-side match,
# and hapax-post-merge-deploy reruns this installer on requester/runtime changes.
Environment=HAPAX_DISPATCH_REDEMPTION_ALLOWED_REQUESTER_SHA256=c3dca72faa466c535f679e66f20ca3d7c5bcf92d5d73404e73b13f01943b40b7
RuntimeDirectory=hapax/coord
RuntimeDirectoryMode=0750
StateDirectory=hapax/dispatch-redemption
StateDirectoryMode=0750
ExecStart=/usr/bin/python3 /home/hapax/.cache/hapax/source-activation/worktree/scripts/hapax-dispatch-redemption-authority --serve --events-path /var/lib/hapax/dispatch-redemption/events.jsonl
Restart=on-failure
RestartSec=2s
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
