# systemd/units/hapax-obs-video50-yuyv-compat-bridge.service
[Unit]
Description=OBS-compatible YUYV bridge from DarkPlaces Xvfb :82 to /dev/video50
After=hapax-darkplaces-v4l2.service
PartOf=hapax-darkplaces-v4l2.service
Conflicts=studio-fx-output.service hapax-obs-video42-yuyv-bridge.service hapax-screwm-meet-camera.service
ConditionPathExists=%h/.config/hapax/enable-darkplaces-runtime
ConditionPathExists=/tmp/.X11-unix/X82
ConditionPathExists=/dev/video50
ConditionPathExists=/usr/bin/ffmpeg
OnFailure=notify-failure@%n.service

[Service]
Type=simple
Restart=on-failure
RestartSec=2s
TimeoutStopSec=3s
KillMode=control-group
SendSIGKILL=yes
ExecStartPre=-/usr/bin/v4l2-ctl -d /dev/video50 -c keep_format=0
ExecStartPre=-/usr/bin/v4l2-ctl -d /dev/video50 --set-fmt-video=width=1920,height=1080,pixelformat=YUYV
ExecStartPre=-/usr/bin/v4l2-ctl -d /dev/video50 --set-parm=30
ExecStartPre=-/usr/bin/v4l2-ctl -d /dev/video50 -c sustain_framerate=1
ExecStartPre=/usr/bin/bash -lc 'for i in {1..40}; do [ -S /tmp/.X11-unix/X82 ] && exit 0; sleep 0.25; done; exit 1'
ExecStart=/usr/bin/ffmpeg -hide_banner -loglevel warning -nostdin -f x11grab -video_size 1920x1080 -framerate 30 -i :82.0+0,0 -vf format=yuyv422 -f v4l2 -pix_fmt yuyv422 -r 30 /dev/video50
