# systemd/units/hapax-cc-task-offer-ready.service
[Unit]
Description=Hapax CC task ready-to-offered promotion reconciler
After=hapax-secrets.service

[Service]
Type=oneshot
WorkingDirectory=%h/.cache/hapax/source-activation/worktree
ExecStart=%h/.cache/hapax/source-activation/worktree/scripts/cc-task-offer-ready --reconcile
Environment=HOME=%h
TimeoutStartSec=60
