# systemd/units/hapax-rebuild-logos.timer
[Unit]
Description=Rebuild Hapax Imagination production binary (5 min check)

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s

[Install]
WantedBy=timers.target
