# systemd/units/hapax-youtube-video-id.service
[Unit]
Description=Hapax YouTube Video-ID Publisher (unblocks chat-monitor)
After=hapax-secrets.service network-online.target
Wants=network-online.target
Requires=hapax-secrets.service
Before=chat-monitor.service
OnFailure=notify-failure@%n.service

[Service]
Type=simple
WorkingDirectory=/home/hapax/projects/hapax-council
ExecStart=%h/.local/bin/uv run python %h/.local/bin/hapax-youtube-video-id-publisher
EnvironmentFile=-/run/user/1000/hapax-secrets.env
Environment=PATH=/home/hapax/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=/home/hapax
Environment=GNUPGHOME=/home/hapax/.gnupg
Environment=PASSWORD_STORE_DIR=/home/hapax/.password-store
Restart=on-failure
RestartSec=10
TimeoutStopSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hapax-youtube-video-id
MemoryMax=256M
CPUQuota=10%

[Install]
WantedBy=default.target
