# systemd/units/hapax-novelty-shift-emitter.timer
[Unit]
Description=Hapax novelty-shift impingement emitter timer (retired; service owns cadence)
ConditionPathExists=/run/hapax/enable-retired-novelty-shift-timer

[Timer]
OnBootSec=2min
AccuracySec=200ms
Unit=hapax-novelty-shift-emitter.service

[Install]
WantedBy=timers.target
