# systemd/units/hapax-quake-live-youtube.service
[Unit]
Description=Hapax Quake in-world YouTube live texture source
After=network-online.target
PartOf=hapax-visual-stack.target
ConditionPathExists=%h/.config/hapax/enable-darkplaces-runtime

[Service]
Type=simple
WorkingDirectory=%h/.cache/hapax/source-activation/worktree
ExecStart=%h/.cache/hapax/source-activation/worktree/scripts/quake-live-media-source.py --source youtube --url-file /dev/shm/hapax-compositor/youtube-video-id.txt --youtube-fallback canary --fps 3 --width 2048 --height 1024 --projection sphere-front --sphere-front-aspect 1.7777777778 --mask none --mask-background 0c0b0d --output /dev/shm/hapax-compositor/quake-live-yt.bgra --meta /dev/shm/hapax-compositor/quake-live-yt.json
Restart=always
RestartSec=3

[Install]
WantedBy=hapax-visual-stack.target
