# systemd/units/hapax-content-resolver.service
[Unit]
Description=Hapax Content Resolver Daemon
After=hapax-secrets.service hapax-imagination-loop.service
Wants=hapax-imagination-loop.service
PartOf=hapax-visual-stack.target
StartLimitIntervalSec=300
StartLimitBurst=5
OnFailure=notify-failure@%n.service

[Service]
Type=simple
EnvironmentFile=/run/user/1000/hapax-secrets.env
WorkingDirectory=%h/.cache/hapax/source-activation/worktree
Environment=PATH=%h/.cache/hapax/source-activation/worktree/.venv/bin:%h/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=%h/.cache/hapax/source-activation/worktree
ExecStartPre=%h/.cache/hapax/source-activation/worktree/scripts/hapax-compositor-runtime-source-check --require-file agents/content_resolver/__main__.py
ExecStart=%h/.cache/hapax/source-activation/worktree/.venv/bin/python -m agents.content_resolver
MemoryMax=512M
Restart=on-failure
RestartSec=10

[Install]
WantedBy=hapax-visual-stack.target
