# systemd/units/hapax-governance-enforcement-event-producer.service
[Unit]
Description=Hapax governance enforcement public-event producer
After=network-online.target hapax-secrets.service
Wants=network-online.target
PartOf=hapax.target

[Service]
Type=simple
WorkingDirectory=%h/.cache/hapax/source-activation/worktree
EnvironmentFile=-/run/user/1000/hapax-secrets.env
Environment=PATH=%h/.cache/hapax/source-activation/worktree/.venv/bin:%h/.local/bin:/usr/local/bin:/usr/bin
Environment=PYTHONPATH=%h/.cache/hapax/source-activation/worktree
Environment=LOG_LEVEL=INFO
Environment=HAPAX_RESEARCH_VEHICLE_PUBLIC_EVENT_PATH=/dev/shm/hapax-public-events/events.jsonl
ExecStart=%h/.cache/hapax/source-activation/worktree/.venv/bin/python -m agents.governance_enforcement_public_event_producer
Restart=on-failure
RestartSec=15
MemoryMax=128M

[Install]
WantedBy=default.target
