# systemd/units/hapax-audio-health-pipewire-xrun.service
[Unit]
Description=Hapax audio health M5 — PipeWire pw-top xrun + buffer-underrun monitor
After=pipewire.service
BindsTo=pipewire.service
PartOf=hapax-audio-health.target
ConditionPathExists=%h/.cache/hapax/rebuild/worktree/pyproject.toml

[Service]
Type=simple
NotifyAccess=all
WorkingDirectory=%h/.cache/hapax/rebuild/worktree
ExecStart=%h/.local/bin/uv run python -m agents.audio_health.m5_pipewire_xrun_daemon
Restart=on-failure
RestartSec=10
WatchdogSec=0
EnvironmentFile=-%h/.config/hapax/audio-health.env
MemoryMax=128M

[Install]
WantedBy=hapax-audio-health.target
