# systemd/units/hapax-pr-admission-auto.timer
[Unit]
Description=Hapax PR admission governor — auto-cycle timer
Documentation=file://%h/projects/hapax-council/scripts/hapax-pr-admission

[Timer]
# Fire before the heavier PR automation, then run on its own staggered cadence.
OnBootSec=45s
OnUnitActiveSec=7min
RandomizedDelaySec=45s
# Persistent so missed boots don't accumulate; AccuracySec slightly loose.
Persistent=true
AccuracySec=30s
Unit=hapax-pr-admission-auto.service

[Install]
WantedBy=timers.target
