# systemd/units/hapax-segment-prep.timer
[Unit]
Description=Hapax daily segment prep timer - Command-R-only prepared artifacts

[Timer]
OnCalendar=*-*-* 04:00:00 UTC
Persistent=true
RandomizedDelaySec=300
Unit=hapax-segment-prep.service

[Install]
WantedBy=timers.target
