# systemd/units/hapax-lufs-panic-cap.service
[Unit]
Description=Hapax broadcast master LUFS-S panic-cap (defense-in-depth)
Documentation=file:///home/hapax/projects/hapax-council/agents/studio_compositor/lufs_panic_cap.py
After=pipewire.service wireplumber.service hapax-secrets.service
BindsTo=pipewire.service wireplumber.service
PartOf=studio-compositor.service
OnFailure=notify-failure@%n.service

[Service]
Type=simple
WorkingDirectory=%h/.cache/hapax/source-activation/worktree
Environment=PATH=%h/.cache/hapax/source-activation/worktree/.venv/bin:%h/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=%h/.cache/hapax/source-activation/worktree
Environment=XDG_RUNTIME_DIR=/run/user/1000
ExecStartPre=%h/.cache/hapax/source-activation/worktree/scripts/hapax-compositor-runtime-source-check --require-file agents/studio_compositor/lufs_panic_cap.py
ExecStart=%h/.cache/hapax/source-activation/worktree/.venv/bin/python -m agents.studio_compositor.lufs_panic_cap
Restart=always
RestartSec=5
# Resource caps: numpy + parec is light; reserve modest memory and a
# single CPU. Prometheus HTTP server on :9484 (pinned in lufs_panic_cap.py).
MemoryMax=512M
CPUWeight=50
# Crash visibility — let journal capture both streams.
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
