# systemd/units/hapax-sprint-tracker.service
[Unit]
Description=Hapax Sprint Tracker — vault-native R&D schedule management
After=hapax-secrets.service

[Service]
Type=oneshot
ExecStart=%h/.cache/hapax/source-activation/worktree/.venv/bin/python -m agents.sprint_tracker
WorkingDirectory=%h/.cache/hapax/source-activation/worktree
Environment=HOME=%h
Environment=PYTHONPATH=%h/.cache/hapax/source-activation/worktree

[Install]
WantedBy=default.target
