# systemd/units/hapax-lane-reaper.service
# Hapax-Parked: true
[Unit]
Description=Hapax Parked legacy tmux lane inventory projection service

[Service]
Type=oneshot
ExecStart=%h/.local/bin/hapax-lane-reaper --threshold 30
Environment=HOME=/home/hapax
TimeoutStartSec=120
