# systemd/units/hapax-rotate-dispatch-recruitment-impingements.timer
# Hapax-Auto-Enable: true
[Unit]
Description=hapax — frequent cap check for dispatch/recruitment/impingement JSONL ledgers

[Timer]
OnBootSec=3min
OnUnitActiveSec=15min
RandomizedDelaySec=60
Unit=hapax-rotate-dispatch-recruitment-impingements.service

[Install]
WantedBy=timers.target
