# systemd/units/hapax-obs-audio-bind.timer
[Unit]
Description=Refresh Hapax OBS PipeWire audio binding

[Timer]
OnBootSec=90s
OnUnitActiveSec=60s
AccuracySec=10s
Unit=hapax-obs-audio-bind.service
Persistent=false

[Install]
WantedBy=timers.target
