# systemd/units/hapax-visual-pool-snapshot-harvester.timer
[Unit]
Description=Harvest compositor snapshots into the Hapax visual pool

[Timer]
Unit=hapax-visual-pool-snapshot-harvester.service
OnBootSec=45s
OnUnitActiveSec=1min
AccuracySec=10s
Persistent=false

[Install]
WantedBy=timers.target
