# systemd/units/hapax-audio-health-topology-drift.service
[Unit]
Description=Hapax audio health M6 — topology drift monitor (PipeWire module signature)
After=pipewire.service wireplumber.service
Requires=pipewire-pulse.service
PartOf=hapax-audio-health.target
ConditionPathExists=%h/.cache/hapax/rebuild/worktree/pyproject.toml

[Service]
Type=notify
WorkingDirectory=%h/.cache/hapax/rebuild/worktree
ExecStart=%h/.local/bin/uv run python -m agents.audio_health.m6_topology_drift_daemon
Restart=on-failure
RestartSec=10
WatchdogSec=600s
EnvironmentFile=-%h/.config/hapax/audio-health.env
MemoryMax=128M

[Install]
WantedBy=hapax-audio-health.target
