# systemd/units/hapax-v4l2-bridge.service
[Unit]
Description=V4L2 Bridge — shmsrc→v4l2sink sidecar for /dev/video42
After=hapax-video42-format-guard.service studio-compositor.service
Requires=hapax-video42-format-guard.service studio-compositor.service
BindsTo=studio-compositor.service
PartOf=studio-compositor.service
ConditionPathExists=%h/.cache/hapax/source-activation/worktree/scripts/hapax-v4l2-bridge

[Service]
Type=simple
WorkingDirectory=%h/.cache/hapax/source-activation/worktree
ExecStartPre=%h/.cache/hapax/source-activation/worktree/scripts/hapax-compositor-runtime-source-check --require-file scripts/hapax-v4l2-bridge
ExecStartPre=%h/.cache/hapax/source-activation/worktree/scripts/hapax-v4l2-video42-format-guard --verify-only
ExecStart=%h/.cache/hapax/source-activation/worktree/scripts/hapax-v4l2-bridge
Restart=on-failure
RestartSec=1s
TimeoutStartSec=45s
Environment=GST_DEBUG=2
Environment=HAPAX_V4L2_BRIDGE_DEVICE=/dev/video42
Environment=HAPAX_V4L2_BRIDGE_ENABLED=1
Environment=HAPAX_V4L2_BRIDGE_SOCKET=/dev/shm/hapax-compositor/v4l2-bridge.sock
Environment=HAPAX_V4L2_BRIDGE_WAIT_SECONDS=60
Environment=HAPAX_COMPOSITOR_SOURCE_ROOT=%h/.cache/hapax/source-activation/worktree

[Install]
WantedBy=default.target
