# systemd/units/hapax-audio-routing-check.service
[Unit]
Description=Hapax audio routing invariant check
After=pipewire.service

[Service]
Type=oneshot
ExecStart=%h/.local/bin/hapax-audio-routing-check
StandardOutput=journal
StandardError=journal
# On failure, log prominently but don't crash anything
SuccessExitStatus=1 2
