# systemd/units/hapax-v4l2-watchdog.timer
[Unit]
Description=v4l2 heartbeat watchdog — poll every 10s for stalled frame output

[Timer]
OnBootSec=60
OnUnitActiveSec=10
AccuracySec=1s

[Install]
WantedBy=timers.target
