# systemd/units/hapax-visual-stack.target
[Unit]
Description=Hapax Visual Stack - production visual surface services
# Production visual pipeline: imagination/reverie render, DMN generates
# content, VLA bridges perception, and studio-compositor feeds OBS/V4L2.
# The Tauri/WebKit hapax-logos preview is decommissioned and must not be
# pulled in by this target.
Wants=hapax-imagination.service hapax-imagination-loop.service hapax-dmn.service hapax-reverie.service hapax-content-resolver.service visual-layer-aggregator.service studio-compositor.service
After=logos-api.service hapax-daimonion.service

[Install]
WantedBy=default.target
