# systemd/units/hapax-preset-bias-heartbeat.service
[Unit]
Description=Hapax preset.bias heartbeat fallback — chain-mutation backstop when LLM recruitment stalls
Documentation=file:///tmp/effect-cam-orchestration-audit-2026-05-02.md
Documentation=file://%h/projects/hapax-council/docs/governance/preset-bias-heartbeat-fallback.md
After=studio-compositor.service hapax-imagination.service
Wants=studio-compositor.service

[Service]
Type=simple
WorkingDirectory=%h/projects/hapax-council
ExecStart=%h/projects/hapax-council/.venv/bin/python -m agents.preset_bias_heartbeat
Restart=always
RestartSec=5s
MemoryHigh=64M
MemoryMax=128M
Environment=CUDA_VISIBLE_DEVICES=
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hapax-preset-bias-heartbeat

[Install]
WantedBy=default.target
