# systemd/units/hapax-audio-ducker.service
[Unit]
Description=Hapax audio ducker — VAD-driven mixer-gain controller
Documentation=file:///home/hapax/projects/hapax-council/docs/superpowers/specs/2026-04-23-livestream-audio-unified-architecture-design.md
After=pipewire.service pipewire-pulse.service wireplumber.service hapax-secrets.service hapax-private-broadcast-leak-guard.service
BindsTo=pipewire.service
Wants=hapax-private-broadcast-leak-guard.service
PartOf=pipewire.service
OnFailure=notify-failure@%n.service

[Service]
Type=simple
WorkingDirectory=%h/.cache/hapax/source-activation/worktree
Environment=HOME=%h
Environment=PATH=%h/.cache/hapax/source-activation/worktree/.venv/bin:%h/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=%h/.cache/hapax/source-activation/worktree
ExecStartPre=%h/.cache/hapax/source-activation/worktree/scripts/hapax-compositor-runtime-source-check --require-file agents/audio_ducker/__main__.py
ExecStart=%h/.cache/hapax/source-activation/worktree/.venv/bin/python -m agents.audio_ducker
Restart=always
RestartSec=2s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hapax-audio-ducker

[Install]
WantedBy=default.target
