# systemd/units/hapax-cc-pr-merge-watcher.service
[Unit]
Description=Hapax CC-task PR-merge auto-close watcher (H9, PR3 of cc-hygiene)
Documentation=file://%h/projects/hapax-council/docs/research/2026-04-26-task-list-hygiene-operator-visibility.md
After=hapax-secrets.service

[Service]
Type=oneshot
ExecStart=%h/.local/bin/uv --directory %h/.cache/hapax/source-activation/worktree run python scripts/cc-pr-merge-watcher.py --repo-root %h/.cache/hapax/source-activation/worktree
WorkingDirectory=%h/.cache/hapax/source-activation/worktree
Environment=HOME=%h
Environment=PYTHONPATH=%h/.cache/hapax/source-activation/worktree
# Killswitch: set HAPAX_CC_HYGIENE_OFF=1 in the environment to skip both
# the H8 PostToolUse hook and this watcher.

[Install]
WantedBy=default.target
