# systemd/units/hapax-lane-reaper.timer
[Unit]
Description=Hapax lane reaper timer — scan every 30 minutes

[Timer]
OnBootSec=10min
OnUnitActiveSec=30min
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
