# systemd/units/hapax-daimonion-quarantine-watchdog.service
[Unit]
Description=Hapax Daimonion private-voice quarantine drift watchdog
After=pipewire.service pipewire-pulse.service wireplumber.service
Wants=pipewire.service pipewire-pulse.service wireplumber.service

[Service]
Type=oneshot
WorkingDirectory=/home/hapax/projects/hapax-council
Environment=PATH=/home/hapax/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=/home/hapax
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=PYTHONUNBUFFERED=1
ExecStart=%h/.local/bin/hapax-daimonion-quarantine-watchdog --enforce --witness-path /dev/shm/hapax-daimonion/quarantine-watchdog.json
SuccessExitStatus=2
TimeoutStartSec=60s
MemoryMax=128M
CPUQuota=10%
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hapax-daimonion-quarantine-watchdog

[Install]
WantedBy=default.target
