# systemd/units/hapax-hkp-export.service
[Unit]
Description=HKP scheduled exporter — project the live SDLC corpus to the cache-only sdlc bundle
After=network.target

[Service]
Type=oneshot
# Cache-only, support-non-authoritative, CPU-only (no GPU): a derived projection, never authority.
ExecStart=/home/hapax/.local/bin/uv run --project /home/hapax/projects/hapax-council python /home/hapax/projects/hapax-council/scripts/hapax-hkp-export-sdlc
Environment=PATH=/home/hapax/.local/bin:/home/hapax/.cargo/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=/home/hapax
Environment=CUDA_VISIBLE_DEVICES=
MemoryMax=512M
CPUQuota=50%
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hapax-hkp-export
