# systemd/units/hapax-ward-atlas-gpu.service
[Unit]
Description=Screwm ward-atlas GPU service — Cairo→GPU renderer (5060 Ti) -> BGRA shm the engine blits
Documentation=file://%h/projects/hapax-council/docs/superpowers/specs/2026-05-30-screwm-cairo-gpu-port-design.md
ConditionPathExists=%h/.config/hapax/enable-darkplaces-runtime
After=hapax-secrets.service hapax-darkplaces-v4l2.service

[Service]
Type=simple
Environment=RUST_LOG=info
# Pin the ward-atlas render to the 5060 Ti; the 5090 stays for the DarkPlaces
# render + TabbyAPI. Footprint <200 MiB (atlas + glyph atlas + RD).
Environment=HAPAX_WARD_ATLAS_GPU=5060
# DORMANT by default: writes the SHADOW path quake-live-ward-atlas.gpu.bgra
# (NOT the engine-read path), so this can run safely for the AVSDLC witness
# without touching the live render. Per-mount cutover lists only the slots this
# service owns; unlisted slots stay shadow and keep their Python producers.
# Environment=HAPAX_WARD_ATLAS_ACTIVE_SLOTS=ward-atlas
Environment=HAPAX_WARD_ATLAS_FPS=2
# Reclaim the ~69%-core the Python Cairo producer pegs; keep off audio-reserved cores.
ExecStart=%h/.local/bin/screwm-ward-atlas
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=default.target
