# systemd/units/hapax-usb-topology-witness.service
[Unit]
Description=Hapax USB topology witness and repair
After=pipewire.service wireplumber.service pipewire-pulse.service
Wants=pipewire.service wireplumber.service pipewire-pulse.service

[Service]
Type=oneshot
Environment=HAPAX_USB_TOPOLOGY_POLICY=%h/.config/hapax/usb-topology-policy.json
ExecStartPre=/usr/bin/udevadm settle --timeout=30
ExecStart=%h/.local/bin/hapax-usb-topology-witness --repair

[Install]
WantedBy=default.target
