# systemd/units/hapax-impingement-sampler.service
[Unit]
Description=Hapax impingement bus sampler — boredom/curiosity rate metrics (ytb-QM2)
Documentation=https://github.com/ryanklee/council/blob/main/agents/quality_observability/impingement_sampler.py
After=network-online.target hapax-secrets.service
Wants=network-online.target
PartOf=hapax.target

[Service]
Type=simple
WorkingDirectory=%h/.cache/hapax/source-activation/worktree
Environment=PYTHONPATH=%h/.cache/hapax/source-activation/worktree
EnvironmentFile=-/run/user/1000/hapax-secrets.env
Environment=PATH=%h/.cache/hapax/source-activation/worktree/.venv/bin:/usr/bin
Environment=HAPAX_IMPINGEMENT_SAMPLER_METRICS_PORT=9496
Environment=HAPAX_IMPINGEMENT_SAMPLER_TICK_S=5
Environment=HAPAX_IMPINGEMENT_NOVELTY_WINDOW_S=300
Environment=HAPAX_LOG_LEVEL=INFO
ExecStart=%h/.cache/hapax/source-activation/worktree/.venv/bin/python -m agents.quality_observability.impingement_sampler
Restart=on-failure
RestartSec=15
MemoryMax=128M

[Install]
WantedBy=default.target
