# systemd/units/hapax-gemini-iota-watchdog.service
[Unit]
Description=Retired Hapax Gemini iota lane watchdog compatibility unit
After=default.target
PartOf=hapax-control-plane.target

[Service]
Type=oneshot
ExecStart=/usr/bin/true
RemainAfterExit=yes

# Retained only so stale installs do not respawn the retired Gemini CLI lane.
# This unit intentionally emits no sd_notify readiness signal; no current repo
# units require or want hapax-gemini-iota-watchdog.service.
# Deployed copies reconcile through the normal post-merge source-activation and
# user-systemd reload path.
MemoryMax=64M
CPUQuota=5%

[Install]
WantedBy=default.target
