# systemd/units/hapax-v4l2-bridge-watchdog.timer
[Unit]
Description=Poll V4L2 bridge egress health

[Timer]
OnBootSec=60
OnUnitActiveSec=10
AccuracySec=1s

[Install]
WantedBy=timers.target
