# systemd/units/hapax-obs-video50-yuyv-compat-bridge.service
[Unit]
Description=OBS-compatible YUYV bridge from DarkPlaces /dev/video52 to /dev/video50
After=hapax-darkplaces-v4l2.service
Wants=hapax-darkplaces-v4l2.service
PartOf=hapax-darkplaces-v4l2.service hapax-visual-stack.target
Conflicts=studio-fx-output.service
ConditionPathExists=%h/.config/hapax/enable-darkplaces-runtime
ConditionPathExists=/dev/video52
ConditionPathExists=/dev/video50
ConditionPathExists=/usr/bin/ffmpeg

[Service]
Type=simple
Restart=always
RestartSec=2s
ExecStartPre=-/usr/bin/v4l2-ctl -d /dev/video50 -c keep_format=0
ExecStartPre=-/usr/bin/v4l2-ctl -d /dev/video50 --set-fmt-video=width=1920,height=1080,pixelformat=YUYV
ExecStart=/usr/bin/ffmpeg -hide_banner -loglevel warning -nostdin -f v4l2 -input_format yuyv422 -video_size 1920x1080 -framerate 60 -i /dev/video52 -vf format=yuyv422 -f v4l2 -pix_fmt yuyv422 /dev/video50

[Install]
WantedBy=hapax-visual-stack.target
