# systemd/units/hapax-darkplaces-bridge.service
[Unit]
Description=Hapax DarkPlaces state bridge — copies live state into game directory
PartOf=hapax-darkplaces.service
After=hapax-darkplaces.service
ConditionPathExists=%h/.config/hapax/enable-darkplaces-runtime

[Service]
Type=simple
ExecStart=%h/projects/hapax-council/scripts/darkplaces-state-bridge.sh
Restart=always
RestartSec=2

[Install]
WantedBy=hapax-darkplaces.service
