# systemd/units/hapax-coord.service
[Unit]
Description=Hapax SBCL+CLOG coordination dashboard (activation worktree)
Documentation=file:/home/hapax/projects/hapax-coord/README.md
ConditionPathExists=%h/.cache/hapax/coord-activation/worktree/scripts/run-dev.sh
ConditionPathExists=%h/.cache/hapax/coord-activation/worktree/.deployed-sha
OnFailure=notify-failure@%n.service

[Service]
Type=simple
Environment=HAPAX_COORD_ROOT=%h/.cache/hapax/coord-activation/worktree
WorkingDirectory=%h/.cache/hapax/coord-activation/worktree
ExecStart=%h/.cache/hapax/coord-activation/worktree/scripts/run-dev.sh --daemon
Restart=on-failure
RestartSec=5s
NoNewPrivileges=true

[Install]
WantedBy=default.target
