# systemd/units/hapax-monitor-stream.service
[Unit]
Description=Operator monitor: livestream program (OBS remap fold) -> mk5 OUT9/AUX8 (left ear)
# Volume-controllable leg: Faderfox fader 7 targets output.hapax-monitor-stream.
After=pipewire.service
Wants=pipewire.service

[Service]
ExecStart=/usr/bin/pw-loopback -n hapax-monitor-stream --capture-props=node.target=hapax-obs-broadcast-remap stream.capture.sink=true audio.channels=2 --playback-props="node.target=alsa_output.usb-MOTU_UltraLite-mk5_UL5LFEC2B0-00.pro-output-0 audio.channels=1 audio.position=[AUX8]"
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
