# systemd/units/hapax-worktree-gc.service
[Unit]
Description=Hapax stale worktree garbage collection
# Doc path resolves through the activated deploy tree, NOT the mutable primary
# worktree this unit deliberately stopped running from (see ExecStart rationale).
# Reinstall + `systemctl --user daemon-reload` after editing this unit, or the
# running timer keeps the old ExecStart/WorkingDirectory until the next deploy.
Documentation=file:%h/.cache/hapax/source-activation/worktree/systemd/README.md

[Service]
Type=oneshot
# Run from the ACTIVATED DEPLOY tree, not the primary worktree: the primary tree is a permanent
# lane pinned on unmerged WIP, so its on-disk script froze at an old version and reported
# removable=0 forever (N-tier deploy-skew). The deploy tree is a registered worktree of the same
# .git (shares object store + sees every worktree/branch) and is always redeployed from main.
ExecStart=%h/.cache/hapax/source-activation/worktree/scripts/hapax-worktree-gc.sh --repo %h/.cache/hapax/source-activation/worktree
WorkingDirectory=%h/.cache/hapax/source-activation/worktree
Environment=HOME=%h
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
