# systemd/units/hapax-audio-health-meta.service
[Unit]
Description=Hapax audio-health suite meta-monitor (M0)
Documentation=file:///home/hapax/projects/hapax-council/docs/runbooks/audio-pipeline-health.md
# Watches each audio-health monitor's Prometheus textfile mtime.
# Emits hapax_audio_health_suite_up and _freshness_seconds gauges.
# Ntfys on monitor-down transition. READ-ONLY: no service restarts.
After=hapax-audio-signal-assertion.service
Wants=hapax-audio-signal-assertion.service
ConditionPathExists=%h/.cache/hapax/rebuild/worktree/pyproject.toml

[Service]
Type=simple
NotifyAccess=all
WorkingDirectory=%h/.cache/hapax/rebuild/worktree
EnvironmentFile=-%h/.config/hapax/audio-health-meta.env
ExecStart=%h/.local/bin/uv run python -m agents.audio_health_meta
Restart=always
RestartSec=15s
WatchdogSec=0
TimeoutStartSec=30s
TimeoutStopSec=10s
Nice=15
Slice=app.slice

[Install]
WantedBy=default.target
