# systemd/units/hapax-lane-reaper.timer
# Hapax-Parked: true
[Unit]
Description=Hapax Parked legacy tmux lane inventory projection timer

[Timer]
OnBootSec=10min
OnUnitActiveSec=30min
RandomizedDelaySec=60
