# systemd/units/hapax-root-failure-intake@.service
[Unit]
# Hapax-Install-Scope: system
Description=Record system-scope Hapax failure for %i
Documentation=file:/home/hapax/.local/lib/hapax-recovery/council/current/manifest.json

[Service]
Type=oneshot
User=hapax
Environment=PATH=/home/hapax/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/local/sbin/hapax-root-failure-intake %i
