# systemd/units/hapax-audio-router.service
[Unit]
Description=Hapax dynamic audio router (Evil Pet + S-4 5 Hz arbiter)
Documentation=https://github.com/ryanklee/hapax-council/blob/main/docs/superpowers/specs/2026-04-21-evilpet-s4-dynamic-dual-processor-design.md
After=hapax-secrets.service pipewire.service
Wants=hapax-secrets.service

[Service]
Type=simple
WorkingDirectory=%h/projects/hapax-council
EnvironmentFile=-%h/.cache/hapax/secrets.env
ExecStart=%h/projects/hapax-council/.venv/bin/python -m agents.audio_router.dynamic_router
Restart=on-failure
RestartSec=2s
# Audio routing is critical-path even when the resident set is normally small:
# keep a soft pressure signal while avoiding a 1G hard ceiling on a 128G host.
MemoryHigh=1G
MemoryMax=2G

# Pure software; no GPU. Isolate from CUDA so a stray import cannot
# acquire VRAM the compositor / TabbyAPI rely on.
Environment=CUDA_VISIBLE_DEVICES=

# Stdout / stderr → systemd journal (default).
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
