# systemd/units/hapax-oom-policy-audit.service
[Unit]
Description=Audit Hapax OOM containment policy drift
OnFailure=notify-failure@%n.service
StartLimitIntervalSec=0

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/hapax-oom-policy-audit --json
SyslogIdentifier=hapax-oom-policy-audit
