# systemd/units/hapax-imagination-watchdog.service
[Unit]
Description=Imagination loop composite liveness watchdog with restart backoff
After=hapax-imagination-loop.service

[Service]
Type=oneshot
ExecStart=%h/.local/bin/hapax-imagination-watchdog.sh
Environment=HAPAX_IMAG_WATCHDOG_FILE=/dev/shm/hapax-imagination/current.json
Environment=HAPAX_IMAG_WATCHDOG_STALE_S=600
Environment="HAPAX_IMAG_WATCHDOG_FRESH_PATHS=/dev/shm/hapax-imagination/uniforms.json /dev/shm/hapax-dmn/observations.json /dev/shm/hapax-sensors/snapshot.json"
Environment=HAPAX_IMAG_WATCHDOG_SOURCE_DIR=/dev/shm/hapax-imagination/sources
Environment=HAPAX_IMAG_WATCHDOG_RESTART_COOLDOWN_S=1800
Environment=HAPAX_IMAG_WATCHDOG_COOLDOWN_S=1800
Environment=HAPAX_IMAG_WATCHDOG_UNIT=hapax-imagination-loop.service
Environment=HOME=/home/hapax
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/1000/bus
SyslogIdentifier=hapax-imagination-watchdog
TimeoutSec=10
