# systemd/units/hapax-review-receipts-export.timer
# Hapax-Auto-Enable: true
[Unit]
Description=Review receipts export timer (the trainyard's verdict feed)

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=false

[Install]
WantedBy=timers.target
