# systemd/units/hapax-music-player.service
[Unit]
Description=Local music player — plays operator-selected tracks via pw-cat (content-source-registry Phase 4a)
Documentation=file:///home/hapax/projects/hapax-council/docs/superpowers/plans/2026-04-23-content-source-registry-plan.md
After=hapax-secrets.service pipewire.service hapax-private-broadcast-leak-guard.service
Wants=pipewire.service hapax-private-broadcast-leak-guard.service
OnFailure=notify-failure@%n.service

[Service]
Type=simple
ExecStart=/home/hapax/.local/bin/uv run --directory /home/hapax/projects/hapax-council python -m agents.local_music_player
Environment=PATH=/home/hapax/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=/home/hapax
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=PYTHONUNBUFFERED=1
Restart=on-failure
RestartSec=5s
MemoryMax=512M
CPUQuota=25%
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hapax-music-player

[Install]
WantedBy=default.target
