# systemd/units/llm-backup.service
[Unit]
Description=Deprecated LLM Stack Backup Compatibility Receipt
After=network-online.target
OnFailure=notify-failure@%n.service

[Service]
Type=oneshot
ExecStart=/bin/bash /home/hapax/projects/hapax-council/systemd/scripts/backup.sh
WorkingDirectory=/home/hapax/projects/hapax-council
Environment=PATH=/home/hapax/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=/home/hapax
MemoryHigh=1G
MemoryMax=2G
CPUQuota=50%
SyslogIdentifier=llm-backup
