# systemd/units/logos-api.service
[Unit]
Description=Logos API — FastAPI backend on :8051
After=network.target llm-stack.service hapax-secrets.service
Wants=llm-stack.service
Requires=hapax-secrets.service
StartLimitIntervalSec=0
OnFailure=notify-failure@%n.service

[Service]
Type=simple
WorkingDirectory=/home/hapax/projects/hapax-council
EnvironmentFile=/run/user/1000/hapax-secrets.env
ExecStart=/home/hapax/.local/bin/uv --directory /home/hapax/projects/hapax-council run python -m logos.api --host 0.0.0.0 --port 8051
Restart=on-failure
RestartSec=5
SuccessExitStatus=SIGTERM
Environment=PATH=/home/hapax/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=/home/hapax
# Critical API surface. Live evidence on 2026-05-16 showed normal
# operation at ~0.93G under the old 1G cap while the host still had
# ~79G RAM available. Keep a hard ceiling for runaway protection, but
# do not let a stale small cap page or kill the API.
MemoryHigh=6G
MemoryMax=8G
SyslogIdentifier=logos-api
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
