# systemd/units/mixer-keepalive.timer
[Unit]
Description=Periodic inaudible pulse to keep L-12 USB path active

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s

[Install]
WantedBy=timers.target
