# systemd/units/mixer-keepalive.service
[Unit]
Description=Send inaudible keepalive pulse to the ZOOM L-12 USB path
After=pipewire.service wireplumber.service
Wants=pipewire.service wireplumber.service

[Service]
Type=oneshot
# 1s 10Hz @ -80dB. This is a presence pulse, not a routing fallback.
Environment=HAPAX_L12_KEEPALIVE_PULSE=%h/.local/share/mixer-keepalive.wav
ExecStart=%h/.cache/hapax/source-activation/worktree/scripts/hapax-l12-keepalive-pulse
MemoryMax=64M
