# systemd/units/notify-failure@.service
[Unit]
Description=Send failure notification for %i

[Service]
Type=oneshot
ExecStart=%h/.cache/hapax/source-activation/worktree/scripts/hapax-p0-incident-intake service-failed %i
