# systemd/units/stimmung-sync.service
[Unit]
Description=Stimmung sync — persist system self-state to RAG
OnFailure=notify-failure@%n.service

[Service]
Type=oneshot
WorkingDirectory=%h/.cache/hapax/source-activation/worktree
ExecStartPre=%h/.cache/hapax/source-activation/worktree/scripts/hapax-compositor-runtime-source-check --require-file agents/stimmung_sync.py
ExecStart=%h/.cache/hapax/source-activation/worktree/.venv/bin/python -m agents.stimmung_sync
Environment=PATH=%h/.cache/hapax/source-activation/worktree/.venv/bin:%h/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=%h/.cache/hapax/source-activation/worktree
Environment=HOME=%h
StandardOutput=journal
StandardError=journal
SyslogIdentifier=stimmung-sync
MemoryHigh=1G
MemoryMax=2G
CPUQuota=20%
