# systemd/units/studio-fx-output.service
[Unit]
Description=Studio FX Output — feeds GPU-effected video to /dev/video50
After=studio-compositor.service
Wants=studio-compositor.service

[Service]
Type=simple
WorkingDirectory=%h/.cache/hapax/source-activation/worktree
ExecStartPre=%h/.cache/hapax/source-activation/worktree/scripts/hapax-compositor-runtime-source-check --require-file systemd/units/studio-fx-output.sh
ExecStart=%h/.cache/hapax/source-activation/worktree/systemd/units/studio-fx-output.sh
Restart=on-failure
RestartSec=5
Environment=PATH=%h/.cache/hapax/source-activation/worktree/.venv/bin:%h/.local/bin:/usr/local/bin:/usr/bin:/bin
SyslogIdentifier=studio-fx-output

[Install]
WantedBy=default.target
