# systemd/units/video-retention.service
[Unit]
Description=Video recording retention — prune old segments

[Service]
Type=oneshot
ExecStart=%h/.local/bin/hapax-video-retention
Environment=HOME=/home/hapax
