# systemd/units/visual-layer-aggregator.service
[Unit]
Description=Visual Layer Aggregator — polls logos API and perception state for compositor overlay
After=logos-api.service hapax-daimonion.service hapax-secrets.service
Requires=hapax-secrets.service
Wants=logos-api.service

[Service]
Type=simple
WorkingDirectory=%h/.cache/hapax/source-activation/worktree
EnvironmentFile=/run/user/1000/hapax-secrets.env
ExecStartPre=%h/.cache/hapax/source-activation/worktree/scripts/hapax-compositor-runtime-source-check --require-file agents/visual_layer_aggregator/__main__.py
ExecStart=%h/.cache/hapax/source-activation/worktree/.venv/bin/python -m agents.visual_layer_aggregator
Restart=on-failure
RestartSec=10
Environment=PATH=%h/.cache/hapax/source-activation/worktree/.venv/bin:%h/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=%h/.cache/hapax/source-activation/worktree
Environment=HOME=%h
MemoryHigh=2G
MemoryMax=4G
SyslogIdentifier=visual-layer-aggregator

[Install]
WantedBy=default.target
