# systemd/units/youtube-player.service
[Unit]
Description=YouTube Player — decode to v4l2loopback + PipeWire audio
After=studio-compositor.service pipewire.service
Wants=studio-compositor.service
StartLimitIntervalSec=5min
StartLimitBurst=6

[Service]
Type=simple
WorkingDirectory=%h/.cache/hapax/source-activation/worktree
Environment=PATH=%h/.cache/hapax/source-activation/worktree/.venv/bin:%h/.local/bin:/usr/local/bin:/usr/bin
Environment=PYTHONPATH=%h/.cache/hapax/source-activation/worktree
Environment=HOME=%h
Environment=XDG_RUNTIME_DIR=%t
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=%t/bus
# YT audio engine only: the YT VISUAL lives in the screwm render's OARB, so the
# redundant v4l2 video loopback (vestige of the retired studio-compositor) is
# disabled — frees /dev/video50 for the screwm Meet-in-screwm camera, no contention.
# Audio (-> hapax-yt-loudnorm -> MPC USB 7/8) + JPEG snapshots are unaffected.
Environment=HAPAX_YOUTUBE_EMIT_V4L2=0
ExecStart=%h/.cache/hapax/source-activation/worktree/.venv/bin/python %h/.cache/hapax/source-activation/worktree/scripts/youtube-player.py
Restart=on-failure
RestartSec=15s
MemoryHigh=1500M
MemoryMax=2G
MemorySwapMax=1G

[Install]
WantedBy=default.target
