# systemd/units/hapax-imagination.service.d/oom-protect.conf
[Service]
# Visual surface producer remains protected after user@1000 is made killable.
OOMScoreAdjust=-800
