# systemd/units/hapax-imagination.service.d/zzzz-screwm-reverie-source-only.conf
[Service]
# Screwm/Quake consumes Reverie as an in-world texture through
# /dev/shm/hapax-sources/reverie.rgba. During the migration this service must
# not also own the livestream V4L2 device, and it must render the Reverie
# substrate rather than the legacy 3D proof scene.
Environment=HAPAX_IMAGINATION_V4L2_OUTPUT=0
Environment=HAPAX_IMAGINATION_3D_PROOF=0
Environment=HAPAX_IMAGINATION_WIDTH=960
Environment=HAPAX_IMAGINATION_HEIGHT=540
Environment=HAPAX_IMAGINATION_INTERVAL_MS=33
Environment=HAPAX_IMAGINATION_OUTPUT_HALF_RATE=1
