# systemd/units/hapax-daimonion.service.d/aec.conf
[Service]
# Read by agents/hapax_daimonion/audio_input.py — prefer the echo-cancelled
# source when the PipeWire AEC loopback module is installed.
Environment=HAPAX_AEC_ACTIVE=1
