# systemd/units/hapax-daimonion.service.d/audio-input.conf
[Service]
# STT capture source — read by agents/hapax_daimonion/audio_input.py. A
# non-empty value overrides DEFAULT_SOURCE_PRIORITY. The ReSpeaker XVF3800
# (Profile 3 = input:analog-stereo) is the sole live STT mic. Replaces the
# unversioned rode-input.conf (Rode Wireless PRO RX — not on the USB bus, a
# dead target) + zz-stale-rode-runtime-mitigation.conf override pair; see
# README.md in this directory.
Environment=HAPAX_AUDIO_INPUT_TARGET=alsa_input.usb-Seeed_Studio_reSpeaker_XVF3800_4-Mic_Array_101991441260500383-00.analog-stereo
