# systemd/units/hapax-daimonion.service.d/capacity.conf
[Service]
# Raises the main unit's 10G-soft/12G-hard ceilings. The working set
# (torch + ctranslate2 + onnxruntime + numba in one process) reached 6.3G
# cgroup-current at 5h uptime; the first-camera vision lazy-load spikes to
# ~7.5G. 12G-soft keeps reclaim pressure off the steady state and 16G-hard
# leaves real margin before a cgroup OOM kill.
MemoryHigh=12G
MemoryMax=16G
MemorySwapMax=4G
