# systemd/units/hapax-daimonion.service.d/cpu-affinity.conf
# Fence daimonion's CPU-side work off the audio data-loop cores {6,7,14,15};
# allow {0-5, 8-13} — the same exclusion the studio-compositor carries.
#
# Why exclude rather than co-locate: daimonion does real, bursty CPU work on the
# host (Kokoro TTS synthesis, vision-fusion lazy-loads of 7+ models, STT
# orchestration). Per systemd/overrides/audio-stability/README.md, the dedicated
# audio cores must stay UNCONTENDED — that is the whole point of pinning the
# SCHED_FIFO 88 data-loops there. A daimonion vision-load spike landing on
# {6,7,14,15} is exactly the kind of contention that caused the 2026-04-27 L-12
# dropouts, so we keep it off, just like the compositor.
#
# daimonion's TTS is still shielded from the SDLC fleet: the unit is Nice=-10 and
# the fleet runs at CPUWeight=idle (SCHED_IDLE) in hapax-sdlc.slice, so the fleet
# yields to daimonion unconditionally even on the shared {0-5,8-13} pool.
#
# (If the operator prefers daimonion CO-LOCATED on the audio island {6,7,14,15}
# instead — TTS next to the data-loops, fully fleet-fenced — flip this to
# `CPUAffinity=6 7 14 15`; it is a one-line change. Exclusion is the audio-safest
# reading and is shipped by default.)
[Service]
CPUAffinity=0 1 2 3 4 5 8 9 10 11 12 13
