# systemd/units/hapax-daimonion.service.d/gpu-pin.conf
[Service]
# Pin daimonion STT to GPU 1 (RTX 5060 Ti, 16GB). GPU 0 (RTX 5090)
# is shared between TabbyAPI (Command-R 35B) and studio compositor.
# faster-whisper uses CTranslate2 — works on sm_120.
Environment="CUDA_DEVICE_ORDER=PCI_BUS_ID"
Environment="CUDA_VISIBLE_DEVICES=1"
