# systemd/units/hapax-daimonion.service.d/gpu-pin.conf
[Service]
# LRR Phase 3 item 1 (corrected) — migrate the hapax-daimonion voice
# stack from GPU 1 (3090) to GPU 0 (5060 Ti) to free the 3090 for
# Hermes 3 70B.
#
# Correction note: the original PR #811 drop-in was placed at
# ``systemd/units/hapax-dmn.service.d/gpu-pin.conf``, targeting the
# wrong service. ``hapax-dmn.service`` runs ``-m agents.dmn`` — a
# CPU-only cognitive substrate (``MemoryMax=1G``, no GPU models).
# The 3.4 GiB GPU footprint beta attributed to "hapax-dmn" is
# actually held by ``hapax-daimonion.service`` (``-m agents.hapax_daimonion``),
# which runs faster-whisper STT + embeddings on CUDA. This drop-in
# fixes the target.
#
# Under Option α (pre-Phase 3), hapax-daimonion's base unit sets no
# CUDA_* env vars. With CUDA default FASTEST_FIRST ordering the
# 3090 is cuda:0, so the voice stack loads there (~3.4 GiB:
# faster-whisper STT + embedding).
#
# Under Option γ, the 3090 must host Hermes 3 70B's primary weight
# slice (~23.5 GiB + activations + KV cache). There is NO room for
# the voice stack on GPU 1 once Hermes 3 is loaded. The daimonion
# moves to GPU 0 (5060 Ti, 16 GiB total) alongside the compositor
# (3.3) and hapax-imagination (0.3). New GPU 0 budget:
#
#   compositor      3.3 GiB
#   imagination     0.3
#   hapax-daimonion 3.4   ← moved here
#   Hermes overflow 2.75  ← Option γ slice
#   faster-whisper  2.5   ← compositor-venv copy, already on GPU 0
#   ─────────────────────
#   total          12.25 GiB   (3.75 GiB headroom on 16 GiB card)
#
# Beta verified the math in the Phase 3 supplement drop (see
# ~/.cache/hapax/relay/context/2026-04-14-beta-phase-3-supplement-
# verified-preconditions.md § 1.4). Beta attributed the 3.4 GiB to
# "hapax-dmn" but the actual systemd unit holding GPU memory is
# hapax-daimonion (the voice daemon). See the `hapax-daimonion`
# module path throughout `agents/hapax_daimonion/` and the
# `presence_engine.py` / `backends/ir_presence.py` imports.
#
# `CUDA_DEVICE_ORDER=PCI_BUS_ID` is REQUIRED so `CUDA_VISIBLE_DEVICES=0`
# resolves to the 5060 Ti and not to whichever card CUDA's default
# FASTEST_FIRST heuristic decides is "fastest". See the companion
# tabbyapi.service.d/gpu-pin.conf drop-in for the same ordering note.
#
# Rollback to Option α: remove this drop-in entirely
# (`rm ~/.config/systemd/user/hapax-daimonion.service.d/gpu-pin.conf`),
# daemon-reload, restart hapax-daimonion. The base unit has no CUDA
# env vars and CUDA default FASTEST_FIRST places the voice stack back
# on the 3090 alongside TabbyAPI.
Environment="CUDA_DEVICE_ORDER=PCI_BUS_ID"
Environment="CUDA_VISIBLE_DEVICES=0"
