# systemd/units/hapax-daimonion.service.d/gpu-pin.conf
[Service]
# Expose both GPUs: Chatterbox TTS → cuda:0 (RTX 5090, 8.5GB free),
# Parakeet STT (ONNX, ~600MB) coexists on whichever GPU ONNX picks.
# Old pin to GPU 1 only caused OOM when Chatterbox + STT + imagination
# exceeded 5060 Ti's 16GB.
Environment="CUDA_DEVICE_ORDER=PCI_BUS_ID"
Environment="CUDA_VISIBLE_DEVICES=0,1"
