# systemd/units/hapax-daimonion.service.d/oom-protect.conf
[Service]
# Voice/TTS path remains protected after user@1000 stops shielding all children.
OOMScoreAdjust=-500
