# systemd/units/hapax-daimonion.service.d/oom-protect.conf
[Service]
# Start at the valid user-manager score, then let the root helper apply -500 to
# the authoritative service cgroup without protecting unrelated user children.
OOMScoreAdjust=100
ExecStartPost=-/usr/local/bin/hapax-oom-score-trigger %n
MemoryLow=2G
MemoryMin=1G
