# systemd/units/hapax-daimonion.service.d/s4-rearm.conf
# A daimonion (re)start can leave the S-4 unarmed: the dynamic_router re-asserts
# the gain ladder only on a scene CHANGE (dynamic_router.py), not on startup, so
# a daimonion-only bounce re-arms nothing and broadcast voice silently dies
# (confirmed live 2026-06-20, 14:46). Pull a silent re-arm on every start. The
# re-arm unit is ordered After=hapax-daimonion.service, so it runs once the
# daimonion/router is up; the recurring timer is the backstop.
[Unit]
Wants=hapax-s4-rearm.service
