# systemd/units/hapax-daimonion.service.d/shutdown-killmode.conf
# Shutdown KillMode hardening — 2026-05-03
#
# Bug 2 RCA: 8 of 27 daimonion process exits today were SEGV core-dumps from
# native worker thread pools outliving Python main thread during shutdown:
#   - libgomp gomp_barrier_wait_end (OpenMP)
#   - mediapipe::ThreadPool::WorkerThread (vision)
#   - tokenizers.abi3.so rayon-core (tokenizers)
#
# KillMode=mixed sends SIGTERM to PID 1 (preserves graceful-shutdown handlers
# in run_inner.py:_main_async finally), then SIGKILL to remaining children
# after TimeoutStopSec. Zero behavior change during normal operation.

# 2026-05-04 update: TimeoutStopSec 5s → 30s. Rebuild-service-driven restarts
# during PR clusters were SIGKILLing daimonion mid-perception-backend-stop
# (95-task python process needs >5s to gracefully stop 16+ perception
# backends + pw-cat children + AEC + workspace_monitor). 30s gives enough
# headroom while still protecting against pathological hangs. The
# KillMode=mixed pattern is unchanged — this only extends the SIGTERM
# grace period before SIGKILL fires.

[Service]
KillMode=mixed
TimeoutStopSec=30s
