# systemd/units/hapax-daimonion.service.d/tts-backend.conf
[Service]
# Kokoro GPU as primary TTS — 40ms synthesis, 312MB VRAM.
# Voice identity lives in S-4 hardware effects, not the TTS model.
# Chatterbox available as fallback for paralinguistic tags when needed.
Environment=HAPAX_TTS_BACKEND=kokoro
Environment=PYTORCH_CUDA_ALLOC_CONF=backend:cudaMallocAsync
