# systemd/units/hapax-daimonion.service.d/tts-server.conf
[Service]
# CASE-VOICE-FOUNDATION-20260610: daimonion is now a TTS client. The
# model-owning process is hapax-tts-local.service, pinned to the podium 5060 Ti,
# so restarting daimonion reloads STT and runtime state but not TTS models.
Environment=HAPAX_TTS_TRANSPORT=server
Environment=HAPAX_TTS_SERVER_SOCKET=/run/user/1000/hapax-tts-local.sock
Environment=HAPAX_TTS_REQUEST_DEADLINE_S=30
