# systemd/units/hapax-daimonion.service.d/tts-backend.conf
[Service]
# TTS primary backend — read by agents/hapax_daimonion/tts.py
# (resolve_backend_from_env). Valid: chatterbox | kokoro. Invalid or unset
# resolves to chatterbox, which falls back to Kokoro when the model cannot
# load (sm_120 has no torch kernels until the torch>=2.9+cu128 keystone,
# Phase 1 of CASE-VOICE-FOUNDATION-20260610). Kokoro synthesizes on CPU
# until that keystone lands. Voice identity lives in S-4 hardware effects,
# not the TTS model.
Environment=HAPAX_TTS_BACKEND=kokoro
# Read by PyTorch, not hapax code.
Environment=PYTORCH_CUDA_ALLOC_CONF=backend:cudaMallocAsync
