# systemd/units/hapax-obs-v4l2-source-reset.service.d/zzzz-screwm-quake-video52.conf
[Service]
ExecStart=
ExecStart=%h/.cache/hapax/source-activation/worktree/.venv/bin/python %h/.cache/hapax/source-activation/worktree/scripts/hapax-obs-v4l2-source-reset --source-name "Video Capture Device (V4L2)" --poll-interval 15 --stall-threshold 60 --reset-cooldown 60 --device-id /dev/video50 --resolution 1920x1080 --framerate 30 --pixelformat YUYV --disable-buffering --auto-reset-input --obs-log-v4l2-errors --ignore-static-screenshot-stalls --max-static-ignore-seconds 75 --screenshot-width 160 --screenshot-height 90 --prom-path %h/.local/share/node_exporter/textfile_collector/hapax-obs-v4l2-source-reset.prom
