# systemd/units/hapax-obs-v4l2-source-reset.service.d/zzzz-screwm-quake-video52.conf
[Service]
ExecStart=
ExecStart=%h/.cache/hapax/source-activation/worktree/.venv/bin/python %h/.cache/hapax/source-activation/worktree/scripts/hapax-obs-v4l2-source-reset --source-name "Video Capture Device (V4L2)" --stall-threshold 60 --reset-cooldown 60 --device-id /dev/video52 --resolution 1920x1080 --framerate 60 --pixelformat YUYV --disable-buffering --auto-reset-input --producer-service hapax-darkplaces-v4l2.service --producer-restart-after-obs-resets 1 --producer-restart-cooldown 300 --producer-restart-timeout 30 --prom-path %h/.local/share/node_exporter/textfile_collector/hapax-obs-v4l2-source-reset.prom
