# systemd/system/dbus-broker.service.d/oom-protect.conf
[Service]
# D-Bus already ships with strong OOM protection; keep that strength explicit
# in the incident drop-in instead of weakening it with the generic value.
OOMScoreAdjust=-900
