# systemd/units/album-identifier.service.d/memory.conf
# Versioned copy of the deployed podium drop-in (was unversioned at
# ~/.config/systemd/user/album-identifier.service.d/memory.conf —
# W4 audit drift finding). PIL decode + duotone colorize on every poll
# tick can balloon; cap before the OOM killer picks a worse victim.
[Service]
MemoryHigh=3G
MemoryMax=4G
MemorySwapMax=2G
