# systemd/units/pipewire-pulse.service.d/oom-protect.conf
[Service]
# Broadcast/private monitor routing is recovery-critical during memory pressure.
OOMScoreAdjust=-900
