# systemd/system/apcupsd.service.d/oom-protect.conf
[Service]
# UPS telemetry must survive memory pressure long enough to distinguish
# real power events from host-local collapse.
OOMScoreAdjust=-900
