# systemd/units/tabbyapi.service.d/gpu-pin.conf
[Service]
# Expose both GPUs to the resident Command-R TabbyAPI process. The
# content-prep path verifies the loaded model before generation; this
# drop-in is only device visibility, not a model swap mechanism.
#
# `CUDA_DEVICE_ORDER=PCI_BUS_ID` keeps CUDA indices aligned with the
# rig's tested gpu_split ordering in tabbyAPI/config.yml. Do not change
# this drop-in as part of segment prep.
Environment="CUDA_DEVICE_ORDER=PCI_BUS_ID"
Environment="CUDA_VISIBLE_DEVICES=0,1"
