# systemd/units/tabbyapi.service.d/gpu-pin.conf
[Service]
# Pin TabbyAPI to GPU 0 (RTX 5090, 32GB). Command-R 35B 5.0bpw fits
# single-GPU with gpu_split=[30]. No split needed post-5090 install.
Environment="CUDA_DEVICE_ORDER=PCI_BUS_ID"
Environment="CUDA_VISIBLE_DEVICES=0"
