# systemd/system/user@1000.service.d/oom.conf
[Service]
# Keep the user manager at the packaged OOM priority. Protect specific recovery
# services with narrow drop-ins instead; do not make every interactive workload
# harder to kill than system login/network/resolver daemons.
OOMScoreAdjust=100
MemoryHigh=96G
MemoryMax=112G
MemorySwapMax=8G
MemoryLow=16G
MemoryMin=8G
