# config/root-required/oom-containment.files
# Versioned ownership manifest for the root-required OOM containment package.
config/root-required/oom-containment.files
scripts/install-p0-oom-containment
config/root-required/hapax-oom-score-enforce.sudoers
systemd/system/system.slice.d/oom-containment.conf
systemd/system/user.slice.d/oom-containment.conf
systemd/system/user-1000.slice.d/oom-containment.conf
systemd/system/user@1000.service.d/oom.conf
systemd/system/apcupsd.service.d/oom-protect.conf
systemd/system/systemd-logind.service.d/oom-protect.conf
systemd/system/systemd-resolved.service.d/oom-protect.conf
systemd/system/systemd-timesyncd.service.d/oom-protect.conf
systemd/system/NetworkManager.service.d/oom-protect.conf
systemd/system/dbus-broker.service.d/oom-protect.conf
systemd/system/sshd.service.d/oom-protect.conf
systemd/units/hapax-root-failure-intake@.service
systemd/units/hapax-oom-score-enforce.service
systemd/units/hapax-oom-score-enforce.timer
systemd/units/hapax-oom-policy-audit.service
systemd/units/hapax-oom-policy-audit.timer
systemd/units/hapax-root-required-deploy-audit.service
systemd/units/hapax-root-required-deploy-audit.timer
systemd/units/app.slice.d/oom-containment.conf
systemd/units/session.slice.d/oom-containment.conf
systemd/units/pipewire.service.d/oom-protect.conf
systemd/units/pipewire-pulse.service.d/oom-protect.conf
systemd/units/wireplumber.service.d/oom-protect.conf
systemd/units/hapax-daimonion.service.d/oom-protect.conf
systemd/units/studio-compositor.service.d/oom-protect.conf
systemd/units/hapax-imagination.service.d/oom-protect.conf
config/earlyoom/default
scripts/hapax-oom-score-enforce
scripts/hapax-oom-score-trigger
scripts/hapax-root-failure-intake
scripts/hapax-oom-policy-audit
scripts/hapax-root-required-deploy-audit
