# systemd/units/wireplumber.service.d/oom-protect.conf
[Service]
# Session graph policy must survive long enough for audio recovery.
OOMScoreAdjust=-900
MemoryLow=512M
MemoryMin=256M
