# systemd/units/wireplumber.service.d/oom-protect.conf
[Service]
# User managers cannot lower child OOM scores. Root applies the protected score
# after startup; the recurring enforcer remains the fail-open repair path.
OOMScoreAdjust=100
ExecStartPost=-/usr/bin/sudo -n /usr/local/sbin/hapax-oom-score-enforce --apply-unit %n
MemoryLow=512M
MemoryMin=256M
