# systemd/units/wireplumber.service.d/oom-protect.conf
[Service]
# WirePlumber carries NoNewPrivileges=yes, so sudo cannot cross the root
# boundary from ExecStartPost. The recurring root enforcer is the repair path.
OOMScoreAdjust=100
MemoryLow=512M
MemoryMin=256M
